# package/lean/luci-app-zerotier/root/etc/uci-defaults/40_luci-zerotier
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@zerotier[-1]
	add ucitrack zerotier
	set ucitrack.@zerotier[-1].init=zerotier
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
