# customfeeds/packages/multimedia/icecast/files/icecast.init
#!/bin/sh /etc/rc.common

# Startup script
# Copyright (C) 2015 OpenWrt.org

USE_PROCD=1
START=90
STOP=15

start_service() {
	user_exists icecast 87 || user_add icecast 87
	group_exists icecast 87 || group_add icecast 87

	[ -d /var/log/icecast ] || {
		mkdir -m 0755 -p /var/log/icecast
		chown icecast:icecast /var/log/icecast
	}

	procd_open_instance
	procd_set_param command /usr/bin/icecast -c /etc/icecast.xml
	procd_set_param respawn
	procd_close_instance
}
