# customfeeds/packages/utils/crelay/files/crelay.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2016 OpenWrt.org

NAME=crelay
START=93

PROG=/usr/bin/${NAME}

USE_PROCD=1

# Custom relay labels (for Web GUI) are defined in /etc/crelay.conf

start_service() {
	procd_open_instance

	# restart if it croaks
	procd_set_param respawn

	procd_set_param command "$PROG"
	# daemon mode in foreground - starts HTTP server
	procd_append_param command -d

	procd_close_instance
}
