# target/linux/lantiq/xway_legacy/base-files/etc/uci-defaults/01_led_migration
#
# Copyright (C) 2013 OpenWrt.org
#

. /lib/functions/migrations.sh

do_internet_led_rename()
{
	config_load system

	[ -n $(config_get led_internet name) ] || return
	[ -z $(config_get led_dsl name) ] || return

	uci rename system.led_internet=led_dsl
	uci set system.led_dsl.name=dsl

	logger -t led-migration "internet led renamed to dsl"
}

case "$(board_name)" in
arcadyan,arv452cqw)
	do_internet_led_rename
	;;
esac

migrations_apply system

exit 0
