# package/lean/luci-app-flowoffload/root/etc/uci-defaults/flowoffload
#!/bin/sh

uci set firewall.@defaults[0].flow_offloading=1
uci commit firewall
exit 0
