# package/lean/luci-app-flowoffload/root/etc/uci-defaults/flowoffload
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@flowoffload[-1]
	add ucitrack flowoffload
	set ucitrack.@flowoffload[-1].init=flowoffload
	commit ucitrack
EOF

/etc/init.d/flowoffload enable

rm -f /tmp/luci-indexcache
exit 0
