# target/linux/ar71xx/base-files/etc/uci-defaults/03_network-vlan-migration
#!/bin/sh
#
# Copyright (C) 2010 OpenWrt.org
#

dev="$(uci -q get network.@switch_vlan[0].device)"
vlan="$(uci -q get network.@switch_vlan[0].vlan)"

if [ "$dev" = "rtl8366s" ] && [ "$vlan" = 0 ]; then
	logger -t vlan-migration "VLAN 0 is invalid for RTL8366s, changing to 1"
	uci set network.@switch_vlan[0].vlan=1
	uci commit network
fi
