# target/linux/ar71xx/base-files/etc/uci-defaults/09_fix-seama-header
#!/bin/sh

. /lib/functions.sh

fix_seama_header() {
	local kernel_size=$(sed -n 's/mtd[0-9]*: \([0-9a-f]*\).*"kernel".*/\1/p' /proc/mtd)

	[ "$kernel_size" ] && mtd -c 0x$kernel_size fixseama firmware
}

board=$(board_name)

case "$board" in
dir-869-a1)
	fix_seama_header
	;;
esac
