# package/lean/luci-app-passwall/root/etc/uci-defaults/luci-app-passwall
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@passwall[-1]
	add ucitrack passwall
	set ucitrack.@passwall[-1].init=passwall
	commit ucitrack
	delete firewall.passwall
	set firewall.passwall=include
	set firewall.passwall.type=script
	set firewall.passwall.path=/var/etc/passwall.include
	set firewall.passwall.reload=1
	commit firewall
EOF

killall -q -9 kcptun haproxy shadowsocks-libev shadowsocksr-libev
/etc/init.d/haproxy disable && /etc/init.d/haproxy stop
rm -rf /tmp/luci-*cache
exit 0
