# package/youridol/luci-app-softethervpn/root/etc/crontabs/root
0 2 * * *  /etc/init.d/softethervpn restart
