# target/linux/rockchip/armv8/base-files/etc/hotplug.d/net/40-net-smp-affinity
#!/bin/sh

[ "$ACTION" = add ] || exit

get_device_irq() {
	local device="$1"
	local line
	local seconds

	# wait up to 10 seconds for the irq/device to appear
	for seconds in $(seq 0 9); do
		line=$(grep -m 1 "${device}\$" /proc/interrupts) && break
		sleep 1
	done
	echo ${line} | sed 's/:.*//'
}

set_interface_core() {
	local core_mask="$1"
	local interface="$2"
	local device="$3"

	[ -z "${device}" ] && device="$interface"

	local irq=$(get_device_irq "$device")

	echo "${core_mask}" > /proc/irq/${irq}/smp_affinity
}

case "$(board_name)" in
friendlyarm,nanopi-r2c|\
friendlyarm,nanopi-r2s|\
xunlong,orangepi-r1-plus|\
embedfire,doornet1)
	set_interface_core 2 "eth0"
	set_interface_core 4 "eth1" "xhci-hcd:usb3"
	;;
friendlyarm,nanopi-r4s|\
huake-cloud,guangmiao-g4c)
	set_interface_core 10 "eth0"
	set_interface_core 20 "eth1"
	;;
esac

