# package/lean/smartdns-install/root/etc/uci-defaults/smartdns-install
#!/bin/sh

nohup /usr/share/smartdns-install/smartdns-install.sh > /dev/null 2>&1 &

exit 0
