# package/community/OpenAppFilter/luci-app-oaf/root/etc/uci-defaults/92_add_user_section
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	set appfilter.user=user
	appfilter.feature=feature
	commit appfilter
EOF
exit 0
