# customfeeds/luci/applications/luci-app-unbound/root/etc/uci-defaults/60_luci-unbound
#!/bin/sh

uci -q batch <<-EOF >/dev/null
  delete ucitrack.@unbound[-1]
  add ucitrack unbound
  set ucitrack.@unbound[-1].init=unbound
  commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0

