# customfeeds/packages/admin/muninlite/files/etc/xinetd.d/muninlite
service munin
{
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= root
	group		= root
	server		= /usr/sbin/muninlite
	disable		= no
}
