# customfeeds/packages/net/dhcp-forwarder/files/dhcp-fwd.init
#!/bin/sh /etc/rc.common

START=50
USE_PROCD=1

LOG_D=/var/log
RUN_D=/var/run
CHR_D=/var/lib/dhcp-fwd
PID_F=$RUN_D/dhcp-fwd.pid

start_service() {
	[ -d $LOG_D ] || mkdir -p $LOG_D
	[ -d $RUN_D ] || mkdir -p $RUN_D
	[ -d $CHR_D ] || mkdir -p $CHR_D && chown dhcp-fwd:dhcp-fwd $CHR_D -R

	procd_open_instance
	procd_set_param command dhcp-fwd -n
	procd_set_param respawn
	procd_close_instance
}

service_triggers()
{
	procd_add_reload_trigger "dhcp"
}
