# target/linux/bcm53xx/base-files/etc/uci-defaults/05_packet_steering
# SPDX-License-Identifier: GPL-2.0-only

uci set network.@globals[0].packet_steering="1"
