# target/linux/bcm53xx/base-files/etc/uci-defaults/09_fix_crc
. /lib/functions.sh

board=$(board_name)

kernel_size=$(sed -n 's/mtd[0-9]*: \([0-9a-f]*\).*"\(kernel\|linux\)".*/\1/p' /proc/mtd)

fixtrx() {
	mtd ${kernel_size:+-c 0x$kernel_size} fixtrx firmware && exit 0
}

fixseama() {
	mtd ${kernel_size:+-c 0x$kernel_size} fixseama firmware && exit 0
}

case "$board" in
dlink,dir-885l | \
dlink,dir-890l)
	fixseama
	;;
*)
	fixtrx
	;;
esac
