# target/linux/bcm53xx/base-files/etc/uci-defaults/09_fix_crc
kernel_size=$(sed -n 's/mtd[0-9]*: \([0-9a-f]*\).*"\(kernel\|linux\)".*/\1/p' /proc/mtd)

mtd ${kernel_size:+-c 0x$kernel_size} fixtrx firmware && exit 0
mtd ${kernel_size:+-c 0x$kernel_size} fixseama firmware && exit 0
exit 1
