# package/lean/luci-app-oscam/root/etc/init.d/epg
#!/bin/sh /etc/rc.common
# copyrights 
# (1)epg dowload Script VER 0.10for openwrt  by zjhzzyf



START=96



start(){


	/usr/bin/epg start
}

stop() {
	echo "not implemented yet!"
}

restart() {
	start
}