# package/qat/firmware/quickassist-c3xxx/files/hotplug.d/qat_dev_processes/01-qat_dev_processes
#!/bin/sh

if [ "$DEVICENAME" != "qat_dev_processes" ] && [ "$ACTION" != "add" ]; then exit 0; fi

chgrp qat /dev/qat_dev_processes
chmod 0660 /dev/qat_dev_processes


