# customfeeds/luci/applications/luci-app-privoxy/root/etc/uci-defaults/40_luci-privoxy
#!/bin/sh

# no longer needed for "Save and Apply" to restart privoxy
# luci-app-privoxy calls /etc/init.d/privoxy reload
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@privoxy[-1]
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache

exit 0
