# target/linux/ath79/tiny/base-files/etc/uci-defaults/04_led_migration
. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
tplink,tl-wr802n-v1)
	migrate_leds "blue:=green:"
	;;
esac

remove_devicename_leds

migrations_apply system

exit 0
