# package/lean/luci-app-gost/root/etc/uci-defaults/gost
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@gost[-1]
	add ucitrack gost
	set ucitrack.@gost[-1].init=gost
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
