# target/linux/ipq40xx/base-files/etc/init.d/adjust_network
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=11
STOP=98

adjust_smp_affinity() {
	test -f /lib/adjust_network.sh && {
		. /lib/adjust_network.sh

		adjust_eth_queue
		adjust_edma_smp_affinity
		adjust_radio_smp_affinity
	}
}

boot() {
	adjust_smp_affinity
}
