# target/linux/ipq40xx/base-files/etc/init.d/mt7615_reboot
#!/bin/sh /etc/rc.common

STOP=60

stop() {
	. /lib/functions.sh

	case $(board_name) in
		hiwifi,c526a)
			rmmod mt7615e
			;;
	esac
}
