# target/linux/ar71xx/base-files/etc/init.d/ar71xx_clock
#!/bin/sh /etc/rc.common
# Copyright (C) 2017 lean <coolsnowwolf@gmail.com>

START=99

start()
{
  a=$(awk 'BEGIN{FS="[ \t]+:[ \t]"} /machine/ {print $2}' /proc/cpuinfo)
  b=$(dmesg | awk -F 'CPU:|MHz' '/Clocks:/{print $2+0}')
  c=${a}" CPU "${b}Mhz
  echo -n $c > /tmp/sysinfo/model
}
