# package/lean/luci-app-music-remote-center/root/etc/uci-defaults/forked-daapd
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@forked-daapd[-1]
	add ucitrack forked-daapd
	set ucitrack.@forked-daapd[-1].init=forked-daapd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
