# package/Twinzo/DevicesComment/root/uci-defaults/luci-devcom
#!/bin/sh

chmod 777 /etc/init.d/devcom 

rm -f /tmp/luci-indexcache
exit 0