# package/lean/luci-app-xunlei/root/etc/xware/xlatest
#!/bin/sh /etc/rc.common

VERSION=`wget -qO- http://dl.lazyzhu.com/file/Thunder/Xware/latest`
sleep 3 && killall wget
xl_configfile="/tmp/etc/xlver"
[ ! -f "$xl_configfile" ] && touch "$xl_configfile"
cat > "$xl_configfile" << EOF
$VERSION
EOF

