# customfeeds/packages/libs/redis/files/redis.init
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=95
STOP=10

REDIS_BIN="/usr/bin/redis-server"
REDIS_CONFIG="/etc/redis.conf"
REDIS_PID="/var/run/redis.pid"
REDIS_DATA="/var/lib/redis"

start_service() {
	mkdir -p "$REDIS_DATA"
	procd_open_instance redis
	procd_set_param command "$REDIS_BIN" "$REDIS_CONFIG"
	procd_close_instance
}
