# package/my/luci-app-flowoffload/root/etc/init.d/flowoffload
#!/bin/sh /etc/rc.common
# Copyright (c) 2011-2015 OpenWrt.org

START=60

CRON_FILE=/etc/crontabs/root
DNSMASQ_RESTART=N
DNS_SERVER="1.0.0.1,1.1.1.1,1.2.4.8,4.2.2.1,4.2.2.2,8.8.4.4,8.8.8.8,9.9.9.9,114.114.114.114,114.114.115.115,119.28.28.28,119.29.29.29,180.76.76.76,210.2.4.8,208.67.220.220,208.67.222.222,223.5.5.5,223.6.6.6,240c::6644,240c::6666,2400:da00::6666,2001:4860:4860::8844,2001:4860:4860::8888,2606:4700:4700::1001,2606:4700:4700::1111"

start_pdnsd() {
	DNS_SERVER=$(uci get flowoffload.@flow[0].dns_server 2>/dev/null)
	[ ! -f /usr/sbin/dnscache ] && ln /usr/bin/dnsforwarder /usr/sbin/dnscache
	DNSCACHE_CONF=/var/run/dnscache/dnscache.conf
	PID_PATH=/var/run/dnscache/
	PID_FILE=dns.pid
	mkdir -p $PID_PATH
	pid=`ps | awk '$5 ~ /\[dnscache\]/ {print $1}'`

	cat > /var/run/dnscache/dnscache.conf <<EOF
LogOn true
LogFileThresholdLength 102400
LogFileFolder /var/run/dnscache
UDPLocal 0.0.0.0:5333
UDPGroup $DNS_SERVER * no
GroupFile
BlockIP 243.185.187.39,46.82.174.68,37.61.54.158,93.46.8.89,59.24.3.173,203.98.7.65,8.7.198.45,78.16.49.15,159.106.121.75,69.63.187.12,31.13.76.8,31.13.64.49
IPSubstituting
BlockNegativeResponse false
Hosts
HostsUpdateInterval 18000
HostsDownloadPath
HostsScript
HostsRetryInterval 30
AppendHosts
BlockIpv6WhenIpv4Exists false
UseCache true
CacheSize 1048576
MemoryCache true
CacheFile
IgnoreTTL false
OverrideTTL -1
MultipleTTL 1
ReloadCache false
OverwriteCache false
DisabledType
DisabledDomain
DisabledList
DomainStatistic false
DomainStatisticTempletFile
StatisticUpdateInterval 29
EOF
	
	dnscache -f $DNSCACHE_CONF -d
	echo dnscache running pid is $pid
	logger -t alex the pid of dnscache is $PID_PATH/$PID_FILE $pid
	echo $pid > $PID_PATH/$PID_FILE

}

stop_pdnsd() {
  kill $(pidof dnscache) >/dev/null 2>&1 || killall -9 dnscache >/dev/null 2>&1 
  echo "Stop DNS Cache"
}

start_bbr(){
		bbr_mod_enable=$(uci get flowoffload.@flow[0].bbr_mod_enable 2>/dev/null)
		echo "bbr enable"
		if [ $bbr_mod_enable -eq 0 ]; then
			lsmod | grep bbr
			if [ "$?" == "0" ]; then
				echo "The BBR module has been loaded."
			else
				insmod bbr
			sleep 3
			fi
			echo "bbr" > /etc/modules.d/bbr
			sed -i '/bbr/d' /etc/sysctl.conf
			echo "net.ipv4.tcp_congestion_control=bbr" >> /etc/sysctl.conf
			sysctl -p
		elif [ $bbr_mod_enable -eq 1 ]; then
			lsmod | grep bbr
			if [ "$?" == "0" ]; then
				rmmod bbr
				sleep 1
				lsmod | grep bbr
				if [ "$?" == "0" ]; then
					echo "The BBR module has been loaded."
				else
					insmod tcp_bbr_mod
					sleep 3
				fi
			else
				insmod tcp_bbr_mod
				sleep 3
			fi
			echo "tcp_bbr_mod" > /etc/modules.d/bbr
			sed -i '/bbr/d' /etc/firewall.user
			echo "net.ipv4.tcp_congestion_control=tcp_bbr_mod" >> /etc/sysctl.conf
			sysctl -p
		fi
}

stop_bbr(){
	echo "disable"
	rmmod bbr
	rmmod tcp_bbr_mod
	echo "net.ipv4.tcp_congestion_control=cubic" >> /etc/sysctl.conf
	sed -i '/bbr/d' /etc/sysctl.conf
	sysctl -p
	sed -i '/cubic/d' /etc/sysctl.conf
	sysctl -p
}

start_fullconenat(){
	nat1_all_enabled=$(uci get flowoffload.@flow[0].nat1_all_enabled 2>/dev/null)
	fullconenat_ip=$(uci get flowoffload.@flow[0].fullconenat_ip 2>/dev/null)
	advanced_setting=$(uci get flowoffload.@flow[0].advanced_setting)
	wan=$(uci get flowoffload.@flow[0].interface)
	upnpd=$(uci get upnpd.config.enabled 2>/dev/null)
	lsmod | grep xt_FULLCONENAT
	if [ "$?" == "0" ]; then
		echo "The fullconenat module has been loaded."
	else
		insmod xt_FULLCONENAT
		sleep 3
	fi
	sed -i '/FULLCONENAT/d' /etc/firewall.user
	sed -i '/zone_wan_postrouting -j MASQUERADE/d' /etc/firewall.user
	if [ $nat1_all_enabled -eq 0 ]; then
		if [ $advanced_setting -eq 0 ]; then			
			echo "iptables -t nat -A zone_wan_prerouting -j FULLCONENAT" >> /etc/firewall.user
			echo "iptables -t nat -A zone_wan_postrouting -s $fullconenat_ip -j FULLCONENAT" >> /etc/firewall.user
			echo "iptables -t nat -A zone_wan_postrouting -j MASQUERADE" >> /etc/firewall.user
		else
			echo "iptables -t nat -A zone_wan_prerouting -i $wan -j FULLCONENAT" >> /etc/firewall.user
			echo "iptables -t nat -A zone_wan_postrouting -o $wan -s $fullconenat_ip -j FULLCONENAT" >> /etc/firewall.user
			echo "iptables -t nat -A zone_wan_postrouting -o $wan -j MASQUERADE" >> /etc/firewall.user
		fi
	elif [ $nat1_all_enabled -eq 1 ]; then
		if [ $advanced_setting -eq 0 ]; then
			echo "iptables -t nat -A zone_wan_prerouting -j FULLCONENAT" >> /etc/firewall.user
			echo "iptables -t nat -A zone_wan_postrouting -j FULLCONENAT" >> /etc/firewall.user
		else
			echo "iptables -t nat -A zone_wan_prerouting -i $wan -j FULLCONENAT" >> /etc/firewall.user
			echo "iptables -t nat -A zone_wan_postrouting -o $wan -j FULLCONENAT" >> /etc/firewall.user
		fi
	fi
	uci set firewall.@zone[1].masq=0
	uci commit firewall

	/etc/init.d/firewall restart
	if [ $upnpd -eq 1 ]; then
	sleep 2
	/etc/init.d/miniupnpd restart
	fi
}

stop_fullconenat(){
	fullconenat_ip=$(uci get flowoffload.@flow[0].fullconenat_ip)
		rmmod xt_FULLCONENAT
		sleep 1
		sed -i '/zone_wan_postrouting -j MASQUERADE/d' /etc/firewall.user
		sed -i '/FULLCONENAT/d' /etc/firewall.user
		uci set firewall.@zone[1].masq=1
		uci commit firewall
		/etc/init.d/firewall restart
}

change_dns() {
 	uci delete dhcp.@dnsmasq[0].server >/dev/null 2>&1
	uci add_list dhcp.@dnsmasq[0].server=127.0.0.1#5333
	uci set dhcp.@dnsmasq[0].resolvfile=/tmp/resolv.conf.auto
	uci set dhcp.@dnsmasq[0].noresolv=0
	uci commit dhcp
}

revert_dns() {
	uci del_list dhcp.@dnsmasq[0].server=127.0.0.1#5333 >/dev/null 2>&1
	uci set dhcp.@dnsmasq[0].resolvfile=/tmp/resolv.conf.auto
	uci set dhcp.@dnsmasq[0].noresolv=0
	uci commit dhcp
}

start(){
	bbr_enable=$(uci get flowoffload.@flow[0].bbr_enabled 2>/dev/null)
	nat1_enable=$(uci get flowoffload.@flow[0].nat1_enabled 2>/dev/null)
	dns=$(uci get flowoffload.@flow[0].dns 2>/dev/null)
	if [ $bbr_enable -eq 1 ];  then
		start_bbr
	fi
	if [ $nat1_enable -eq 1 ];  then
		start_fullconenat
	fi
	if [ $dns -eq 1 ];  then
		start_pdnsd
		change_dns
	fi
	uci set firewall.@defaults[0].flow_offloading=$(uci get flowoffload.@flow[0].flow_offloading)
	uci set firewall.@defaults[0].flow_offloading_hw=$(uci get flowoffload.@flow[0].flow_offloading_hw)
	uci commit firewall
	if [ "$DNSMASQ_RESTART" = N ]; then
		/etc/init.d/dnsmasq restart && echo "DNSMASQ change"
		/etc/init.d/firewall restart
	fi
}

stop(){
	bbr_enable=$(uci get flowoffload.@flow[0].bbr_enabled 2>/dev/null)
	nat1_enable=$(uci get flowoffload.@flow[0].nat1_enabled 2>/dev/null)
    dns=$(uci get firewall.@defaults[0].dns 2>/dev/null)
	stop_bbr
	stop_fullconenat
	stop_pdnsd
	revert_dns
	uci set firewall.@defaults[0].flow_offloading=$(uci get flowoffload.@flow[0].flow_offloading)
	uci set firewall.@defaults[0].flow_offloading_hw=$(uci get flowoffload.@flow[0].flow_offloading_hw)
	uci commit firewall
		if [ "$DNSMASQ_RESTART" = N ]; then
		/etc/init.d/dnsmasq restart && echo "DNSMASQ revert"
		/etc/init.d/firewall restart
	fi
}

restart(){
    DNSMASQ_RESTART=Y
    stop
    start
	  /etc/init.d/dnsmasq restart && echo "DNSMASQ restart"
    /etc/init.d/firewall restart
}



