# package/lean/luci-app-brook-pro/root/etc/uci-defaults/brook
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete firewall.brook
	set firewall.brook=include
	set firewall.brook.type=script
	set firewall.brook.path=/etc/brook.include
	set firewall.brook.reload=1
	commit firewall
EOF

/etc/init.d/brookpro stop
/etc/init.d/brookpro disable

rm -f /tmp/luci-indexcache
exit 0
