# package/diy/luci-app-hypermodem/root/etc/uci-defaults/luci-hypermodem
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@hypermodem[-1]
	add ucitrack hypermodem
	set ucitrack.@hypermodem[-1].init=hypermodem
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
