# target/linux/ath79/nand/base-files/etc/hotplug.d/firmware/11-ath10k-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"ath10k/cal-pci-0000:00:00.0.bin")
	case $board in
	domywifi,dw33d)
		caldata_extract "art" 0x5000 0x844
		ath10k_patch_mac $(mtd_get_mac_binary art 0x12)
		;;
	glinet,gl-ar750|\
	glinet,gl-ar750s-nor|\
	glinet,gl-ar750s-nor-nand)
		caldata_extract "art" 0x5000 0x844
		ath10k_patch_mac $(macaddr_add $(mtd_get_mac_binary art 0x0) +1)
		;;
	glinet,gl-x750-nor|\
	glinet,gl-x750-nor-nand)
		caldata_extract "art" 0x5000 0x844
		ath10k_patch_mac $(macaddr_add $(mtd_get_mac_binary art 0x0) +3)
		;;
	zyxel,nbg6716)
		caldata_extract "art" 0x5000 0x844
		ath10k_patch_mac $(macaddr_add $(mtd_get_mac_ascii u-boot-env ethaddr) +1)
		;;
	esac
	;;
"ath10k/pre-cal-pci-0000:00:00.0.bin")
	case $board in
	nec,wg1200cr)
		ath10kcal_extract "art" 20480 12064
		ath10kcal_patch_mac_crc $(mtd_get_mac_ascii devdata wlan5mac)
		ln -sf /lib/firmware/ath10k/pre-cal-pci-0000\:00\:00.0.bin \
			/lib/firmware/ath10k/QCA9888/hw2.0/board.bin
		;;
	netgear,ex6400|\
	netgear,ex7300)
		ath10kcal_extract "caldata" 20480 12064
		ath10kcal_patch_mac $(mtd_get_mac_binary caldata 12)
		;;
	glinet,gl-x1200|\
	glinet,gl-x1200-nor|\
	glinet,gl-x1200-nor-nand)
		ath10kcal_extract "art" 20480 12064
		ln -sf /lib/firmware/ath10k/pre-cal-pci-0000\:00\:00.0.bin \
			/lib/firmware/ath10k/QCA9888/hw2.0/board.bin
		#ath10kcal_patch_mac $(macaddr_add $(cat /sys/class/net/eth0/address) +1)
		;;
	phicomm,k2t)
		ath10kcal_extract "art" 20480 12064
		ath10kcal_patch_mac_crc $(k2t_get_mac "5g_mac")
		ln -sf /lib/firmware/ath10k/pre-cal-pci-0000\:00\:00.0.bin \
			/lib/firmware/ath10k/QCA9888/hw2.0/board.bin
		;;

	esac
	;;
*)
	exit 1
	;;
esac
