# customfeeds/packages/net/atftp/files/atftpd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2020 OpenWrt.org

START=95
PIDFILE=/tmp/run/atftpd.pid

start() {
	config_load atftpd
	config_get SRV service path "/srv/tftp"
	config_get PORT service port 69

	atftpd --pidfile $PIDFILE --user root.root --port $PORT --daemon $SRV
}

stop() {
	kill $(cat $PIDFILE)
}
