# customfeeds/packages/utils/lvm2/files/lvm2.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009 Stefan Monnier
START=15

start () {
	mkdir -p /tmp/lvm/cache
	/sbin/lvm vgscan --ignorelockingfailure --mknodes || :
	/sbin/lvm vgchange -aly --ignorelockingfailure || return 2
}

stop () {
	mkdir -p /tmp/lvm/cache
	/sbin/lvm vgchange -aln --ignorelockingfailure || return 2
}
