# package/lienol/luci-app-control-timewol/root/etc/uci-defaults/luci-app-control-timewol
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@timewol[-1]
	add ucitrack timewol
	set ucitrack.@timewol[-1].init=timewol
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
