# target/linux/rockchip/armv8/base-files/etc/uci-defaults/12_enable-netifd-smp-tune
#!/bin/sh
uci -q batch <<-EOF >/dev/null
	set network.globals.packet_steering=1
	commit network
EOF

exit 0
