# target/linux/rockchip/armv8/base-files/etc/uci-defaults/yyy-usb-wifi-up
#!/bin/sh

sleep 10
[ -n "$(cat /proc/net/wireless | grep wlan)" ] && wifi config && wifi