# package/qat/firmware/quickassist-c3xxx/files/hotplug.d/usdm_drv/01-usdm
#!/bin/sh

if [ "$DEVICENAME" != "usdm_drv" ] && [ "$ACTION" != "add" ]; then exit 0; fi

chgrp qat /dev/usdm_drv
chmod 0660 /dev/usdm_drv


