# customfeeds/packages/utils/podman/files/podman.init
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1
NAME=podman
PROG=/usr/bin/podman

start_service() {
	[ -d /run/podman ] || mkdir -p /run/podman
	procd_open_instance "podman"
	procd_set_param command "$PROG" system service --time=0
	procd_append_param command $command
	procd_close_instance
}

shutdown() {
	${PROG} pod stop -a
	${PROG} container stop -a
}
