# package/kernel/broadcom-wl/files/etc/hotplug.d/net/00-broadcom-wifi-detect
#!/bin/sh

[ "${ACTION}" = "add" ] && [ "${INTERFACE%%[0-9]}" = "wl" ] && {
	/sbin/wifi config
}
