# package/lean/dvb/oscam/files/oscam.init
#!/bin/sh /etc/rc.common

START=99
OSCAM_PATH=/usr/bin

[ -f  /tmp/oscam ] && 
{
	OSCAM_PATH=/tmp
	[ -x  /tmp/oscam ] ||  chmod +x /tmp/oscam
}
	
start() {
	logger "Starting Oscam"
	[ -d /var/log/oscam ] || mkdir -p /var/log/oscam
	service_start ${OSCAM_PATH}/oscam -b -r 2 -u
}

stop() {
	logger "Stopping Oscam"
	service_stop ${OSCAM_PATH}/oscam -b -r 2 -u
}


