# target/linux/ramips/base-files/etc/hotplug.d/usb/10-motion
[ "$ACTION" = "motion" ] && logger webcam motion event
