# package/network/utils/linkwatch/files/linkwatch.init
#!/bin/sh /etc/rc.common

START=99
STOP=99
LINKWATCH_BIN="/usr/sbin/linkwatch"

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

start() {
        service_start $LINKWATCH_BIN
}

stop() {
        service_stop $LINKWATCH_BIN
}

reload() {
        service_reload $LINKWATCH_BIN
}
