# package/lean/luci-app-mosdns/root/etc/uci-defaults/luci-mosdns
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@mosdns[-1]
	add ucitrack mosdns
	set ucitrack.@mosdns[-1].init=mosdns
	commit ucitrack
EOF

rm -rf /tmp/luci-*
exit 0
