# package/diy/3ginfo/files-text/etc/uci-defaults/99-3ginfo
#!/bin/sh

ln -s /usr/share/3ginfo/cgi-bin/3ginfo.sh /usr/bin/3ginfo
exit 0
