# customfeeds/luci/applications/luci-app-statistics/root/etc/uci-defaults/40_luci-statistics
#!/bin/sh

# register commit handler
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@luci_statistics[-1]
	add ucitrack luci_statistics
	set ucitrack.@luci_statistics[-1].init=luci_statistics
	commit ucitrack
EOF

# symlink for busybox httpd
[ -x /usr/sbin/httpd ] && [ ! -h /www/rrdimg ] && \
	ln -s /tmp/rrdimg /www/rrdimg

rm -f /tmp/luci-indexcache
exit 0
