# customfeeds/packages/utils/apparmor/files/apparmor.init
#!/bin/sh /etc/rc.common

START=75
USE_PROCD=1

. /lib/functions/apparmor.sh

restart() {
	apparmor_restart
}

start_service() {
	apparmor_start
}

stop_service() {
	apparmor_stop
}

reload_service() {
	apparmor_reload
}
