# package/lean/luci-app-passwall/root/etc/init.d/passwall
#!/bin/sh /etc/rc.common

START=99
STOP=15

boot() {
	/usr/share/passwall/app.sh boot
}

start() {
	/usr/share/passwall/app.sh start
}

stop() {
	/usr/share/passwall/app.sh stop
}

restart() {
	stop
	start
}