# package/qat/firmware/quickassist-c3xxx/files/hotplug.d/qat_adf_ctl/01-qat_adf_ctl
#!/bin/sh

if [ "$DEVICENAME" != "qat_adf_ctl" ] && [ "$ACTION" != "add" ]; then exit 0; fi

chgrp qat /dev/qat_adf_ctl
chmod 0660 /dev/qat_adf_ctl


