# customfeeds/packages/net/lcdringer/files/lcdringer.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009-2012 OpenWrt.org

START=60

start() {
	/usr/bin/lcdringer &
	echo $! > /var/run/lcdringer.pid
}

stop() {
	if [ ! -f /var/run/lcdringer.pid ]; then
		return
	fi
	/bin/kill -9 `cat /var/run/lcdringer.pid`
}
