# target/linux/ramips/mt7621/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	alfa-network,quad-e4g)
		[ -n "$(fw_printenv bootcount bootchanged 2>/dev/null)" ] &&\
			echo -e "bootcount\nbootchanged\n" | /usr/sbin/fw_setenv -s -
		;;
	linksys,ea7500-v2)
		mtd resetbc s_env || true
		;;
	samknows,whitebox-v8)
		fw_setenv bootcount 0
		;;
	esac
}
