# target/linux/ramips/mt7621/base-files/etc/init.d/modem-enable
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	huasifei,ws1208v2)
		echo 1 > /sys/class/gpio/modem_reset/value
		;;
	esac
}
