# target/linux/ramips/mt7621/base-files/etc/init.d/set-model-name
#!/bin/sh /etc/rc.common

START=99

start() {
	[ -n "$(fw_printenv -n model_name 2>/dev/null)" ] && \
		fw_printenv -n model_name > /tmp/sysinfo/model
}
