# customfeeds/packages/utils/pcmciautils/files/etc/hotplug.d/pcmcia_socket/pcmcia_socket.agent
#!/bin/sh
# Copyright (C) 2008 OpenWrt.org

case "$ACTION" in
	add)
		echo PCMCIA: socket number: $SOCKET_NO

		sleep 1;
		/usr/sbin/pcmcia-socket-startup $SOCKET_NO
		DONE_FILE=/sys/class/pcmcia_socket/pcmcia_socket$SOCKET_NO/available_resources_setup_done
		if [ -f  $DONE_FILE ]; then
			if [ -w $DONE_FILE ]; then
				echo -n "42" > $DONE_FILE
			fi;
		fi;
		;;

esac
