# customfeeds/packages/net/apcupsd/files/apcupsd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50
STOP=50

start() {
	/usr/sbin/apcupsd -f /etc/apcupsd/apcupsd.conf
}

stop() {
	kill $(cat /var/run/apcupsd.pid)
}
