# target/linux/omap24xx/base-files/etc/hotplug.d/firmware/10-bme-pmm-image
#!/bin/sh

n810_cal_bme_pmm_extract() {
	[ -x /usr/bin/calvaria ] && {
		/usr/bin/calvaria -p -n bme -i last /dev/mtdblock1 >/lib/firmware/n810-cal-bme-pmm.fw ||\
			echo "CAL-BME extract: Failed to extract blob"
	}
}

. /lib/functions.sh

[ "$FIRMWARE" = "n810-cal-bme-pmm.fw" ] && {
	[ "$(board_name)" != "nokia,n810" ] || {
		[ -e /lib/firmware/n810-cal-bme-pmm.fw ] ||\
			n810_cal_bme_pmm_extract
	}
}
