# customfeeds/luci/applications/luci-app-ahcp/root/etc/uci-defaults/40_luci-ahcp
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@ahcpd[-1]
	add ucitrack ahcpd
	set ucitrack.@ahcpd[-1].init=ahcpd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
