# customfeeds/packages/libs/unixodbc/files/odbc.init
#!/bin/sh /etc/rc.common

START=50

gen_odbcinst() {
	local inifile

	echo "[ODBC]"
	echo "Trace = off"
	echo "TraceFile ="

	for inifile in /etc/odbcinst.ini.d/*.ini; do
		cat "$inifile"
	done
}

start() {
	[ ! -d /tmp/etc ] && mkdir /tmp/etc

	gen_odbcinst > /tmp/etc/odbcinst.ini.new
	chmod 0644 /tmp/etc/odbcinst.ini.new

	[ -e /tmp/etc/odbcinst.ini ] && ( rm /tmp/etc/odbcinst.ini || return 1 )

	mv /tmp/etc/odbcinst.ini.new /tmp/etc/odbcinst.ini
}
