# package/my/luci-app-oscam/root/etc/uci-defaults/luci-oscam
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@oscam[-1]
	add ucitrack oscam
	set ucitrack.@oscam[-1].init=oscam
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
