# customfeeds/packages/utils/open-vm-tools/files/vmtoolsd.init
#!/bin/sh /etc/rc.common

START=10

USE_PROCD=1
PROG=/bin/vmtoolsd

start_service() {
	procd_open_instance
	procd_set_param command "$PROG"
	procd_set_param respawn
	procd_close_instance
}
