# target/linux/ramips/base-files/etc/uci-defaults/04_led_migration
. /lib/functions.sh
. /lib/functions/migrations.sh

remove_devicename_leds "rt2800soc-phy0" "rt2800pci-phy0"

migrations_apply system

exit 0
