# target/linux/ramips/base-files/etc/uci-defaults/09_fix-checksum
#!/bin/sh
#
# Copyright (C) 2012 OpenWrt.org
#

. /lib/functions.sh

fix_checksum() {
	local kernel_size=$(sed -n 's/mtd[0-9]*: \([0-9a-f]*\).*"kernel".*/\1/p' /proc/mtd)

	[ "$kernel_size" ] && mtd -c 0x$kernel_size fix$1 firmware
}

board=$(board_name)

case "$board" in
cy-swr1100 | \
dch-m225 | \
dir-645 | \
dir-860l-b1)
	fix_checksum seama
	;;
dlink,dap-1522-a1)
	fix_checksum wrg
	;;
esac
