# package/lean/luci-app-transparent-proxy/files/root/etc/uci-defaults/luci-transparent-proxy
#!/bin/sh
uci get transparent-proxy.@general[-1] >/dev/null 2>&1 || \
	uci add transparent-proxy general >/dev/null 2>&1
uci get transparent-proxy.@access_control[-1] >/dev/null 2>&1 || \
	uci add transparent-proxy access_control >/dev/null 2>&1
uci commit transparent-proxy
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@transparent-proxy[-1]
	add ucitrack transparent-proxy
	set ucitrack.@transparent-proxy[-1].init=transparent-proxy
	commit ucitrack
	delete firewall.transparent_proxy
	set firewall.transparent_proxy=include
	set firewall.transparent_proxy.type=script
	set firewall.transparent_proxy.path=/var/etc/transparent-proxy.include
	set firewall.transparent_proxy.reload=1
	commit firewall
EOF
exit 0
