# package/Twinzo/luci-app-chinadns/root/usr/lib/chinadns/update
#!/bin/sh /etc/rc.common
# by rapistor

start() {
	curl -k 'http://ftp.apnic.net/apnic/stats/apnic/delegated-apnic-latest' | awk -F\| '/CN\|ipv4/ { printf("%s/%d\n", $4, 32-log($5)/log(2)) }'  > /tmp/chinadns_chnroute.txt
	if [ `grep '1.0.1.0/24' /tmp/chinadns_chnroute.txt` ]; then
		mv /tmp/chinadns_chnroute.txt /etc/chinadns_chnroute.txt
		grep '100.100.100.100' /etc/chinadns_chnroute.txt || echo '100.100.100.100' >>/etc/chinadns_chnroute.txt
		/etc/init.d/chinadns stop
		/etc/init.d/chinadns start
	else
		rm -rf /tmp/chinadns_chnroute.txt
	fi
}
