# target/linux/ipq806x/base-files/etc/uci-defaults/05_fix-compat-version
. /lib/functions.sh

case "$(board_name)" in
	linksys,ea7500-v1|\
	linksys,ea8500)
		uci set system.@system[0].compat_version="2.0"
		uci commit system
	;;
esac

exit 0
