# package/lean/luci-app-docker/root/etc/uci-defaults/luci-docker
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@dockerd[-1]
	add ucitrack dockerd
	set ucitrack.@dockerd[-1].init=dockerd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
