# customfeeds/packages/utils/lcd4linux/files/lcd4linux.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007-2015 OpenWrt.org

START=98
USE_PROCD=1

PROG="/usr/bin/lcd4linux"
CONFIGFILE="/etc/lcd4linux.conf"

start_service() {
	[ -f /etc/lcd4linux.conf ] || return 1
	procd_open_instance
	procd_set_param command "$PROG"
	procd_append_param command -F
	procd_append_param command -q
	procd_append_param command -o /tmp/lcd4linux.png
	procd_append_param command -f "$CONFIGFILE"
	procd_set_param file "$CONFIGFILE"
	procd_close_instance
}
