# customfeeds/luci/applications/luci-app-vpnbypass/root/etc/uci-defaults/40_luci-vpnbypass
#!/bin/sh
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@vpnbypass[-1]
	add ucitrack vpnbypass
	set ucitrack.@vpnbypass[-1].init=vpnbypass
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0

