# target/linux/ath79/nand/base-files/etc/hotplug.d/ieee80211/10-fix-wifi-mac
[ "$ACTION" = "add" ] || exit 0

PHYNBR=${DEVPATH##*/phy}

[ -n $PHYNBR ] || exit 0

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)

case $board in
	zyxel,nbg6716)
		# Set mac address for 2.4g device
		[ "$PHYNBR" -eq 1 ] && \
			mtd_get_mac_ascii u-boot-env ethaddr > /sys${DEVPATH}/macaddress
		;;
esac
