# package/lean/luci-app-n2n_v2/root/etc/uci-defaults/luci-n2n_v2
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@n2n_v2[-1]
	add ucitrack n2n_v2
	set ucitrack.@n2n_v2[-1].init=n2n_v2
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
