# package/base-files/files/etc/uci-defaults/10_migrate-shadow
#!/bin/sh

ppwd="$(sed -ne '/^root:/s/^root:\([^:]*\):.*$/\1/p' /etc/passwd)"
spwd="$(sed -ne '/^root:/s/^root:\([^:]*\):.*$/\1/p' /etc/shadow)"

if [ -n "${ppwd#[\!x]}" ] && [ -z "${spwd#[\!x]}" ]; then
	logger -t migrate-shadow "Moving root password hash into shadow database"
	sed -i -e "s:^root\:[^\:]*\::root\:x\::"     /etc/passwd
	sed -i -e "s:^root\:[^\:]*\::root\:$ppwd\::" /etc/shadow
fi

exit 0
