# package/base-files/files/etc/uci-defaults/13_fix-group-user
. /lib/functions.sh

# Skip if we don't have /usr/lib/opkg/info (APK installation)
[ -d /usr/lib/opkg/info ] || exit 0

for file in $(grep -sl Require-User /usr/lib/opkg/info/*.control); do
	file="${file##*/}"
	file="${file%.control}"
	add_group_and_user "${file}"
done

exit 0
