# package/base-files/files/etc/uci-defaults/14_network-generate-clientid
[ "$(uci -q get network.globals.dhcp_default_duid)" != "auto" ] && exit 0

uci -q batch <<-EOF >/dev/null
	# DUID-UUID - RFC6355
	set network.globals.dhcp_default_duid="$(hexdump -vn 16 -e '"0004" 2/2 "%x"' /dev/urandom)"
	commit network
EOF

exit 0
