# package/network/config/qos-scripts/files/etc/hotplug.d/iface/10-qos
#!/bin/sh
[ "$ACTION" = ifup ] && /etc/init.d/qos enabled && /usr/lib/qos/generate.sh interface "$INTERFACE" | sh
