# package/lean/luci-app-guest-wifi/root/etc/uci-defaults/luci-app-guest-wifi
#!/bin/sh 

uci -q batch <<-EOF >/dev/null
          delete ucitrack.@guest-wifi[-1]
          add ucitrack guest-wifi
          set ucitrack.@guest-wifi[-1].init=guest-wifi
          commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
