# package/lean/luci-app-oscam/root/etc/uci-defaults/luci-sattv
#!/bin/sh 

uci -q batch <<-EOF >/dev/null
          delete ucitrack.@epg[-1]
          add ucitrack epg
          set ucitrack.@epg[-1].init=epg
          delete ucitrack.@oscam[-1]
          add ucitrack oscam
          set ucitrack.@oscam[-1].init=oscam
          add_list ucitrack.@oscam[-1].affects=cron    
          commit ucitrack
EOF
rm -f /tmp/luci-indexcache
exit 0
