# luci-theme-kucat/root/etc/uci-defaults/30_luci-kuacat
#!/bin/sh
uci batch <<-EOF
	set luci.themes.KuCat=/luci-static/kucat
	set luci.main.mediaurlbase=/luci-static/kucat
	commit luci
EOF
rm -Rf /var/luci-modulecache
rm -Rf /var/luci-indexcache
exit 0
