# package/base-files/files/etc/hotplug.d/net/00-sysctl
#!/bin/sh

if [ "$ACTION" = add ]; then
	for CONF in /etc/sysctl.conf /etc/sysctl.d/*.conf; do
		[ ! -f "$CONF" ] && continue;
		sed -ne "/^[[:space:]]*net\..*\.$DEVICENAME\./p" "$CONF" | \
			sysctl -e -p - | logger -t sysctl
	done
fi
