# customfeeds/packages/utils/dbus/files/dbus.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007-2011 OpenWrt.org

START=60

SERVICE_PID_FILE=/var/run/dbus.pid

start() {
	mkdir -m 0755 -p /var/lib/dbus
	mkdir -m 0755 -p /var/run/dbus
	[ -x /usr/bin/dbus-uuidgen ] && /usr/bin/dbus-uuidgen --ensure
	service_start /usr/bin/dbus-daemon --system
}

stop() {
	service_stop /usr/bin/dbus-daemon
}
