# package/Twinzo/luci-app-dogcom/root/etc/uci-defaults/luci-dogcom
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@dogcom[-1]
	add ucitrack dogcom
	set ucitrack.@dogcom[-1].init=dogcom
	commit ucitrack
EOF

chmod 777 /etc/init.d/dogcom
chmod 777 /etc/init.d/dogcom_log
chmod 777 /etc/init.d/dogcom_proxy

rm -f /tmp/luci-indexcache
exit 0
