# target/linux/bcm47xx/base-files/etc/uci-defaults/03_network_migration
#
# Copyright (C) 2014-2015 OpenWrt.org
#

uci show network | grep "\.vlan=0"
[ $? -ne 0 ] && exit 0

logger -t network "network config is invalid, creating new one"

lan_proto="$(uci -q get network.lan.proto)"
lan_ipaddr="$(uci -q get network.lan.ipaddr)"
lan_netmask="$(uci -q get network.lan.netmask)"
wan_proto="$(uci -q get network.wan.proto)"
wan_ipaddr="$(uci -q get network.wan.ipaddr)"
wan_netmask="$(uci -q get network.wan.netmask)"

echo "" > /etc/config/network
config_generate

uci set network.lan.proto=$lan_proto
uci set network.lan.ipaddr=$lan_ipaddr
uci set network.lan.netmask=$lan_netmask
uci set network.wan.proto=$wan_proto
uci set network.wan.ipaddr=$wan_ipaddr
uci set network.wan.netmask=$wan_netmask
uci commit network

exit 0
