# package/lean/luci-app-ssr-pro/root/etc/uci-defaults/ssrpro
#!/bin/sh

uci -q batch <<-EOF >/dev/null
  delete ucitrack.@ssrpro[-1]
	add ucitrack ssrpro
	set ucitrack.@ssrpro[-1].init=ssrpro
	commit ucitrack
	delete firewall.ssrpro
	set firewall.ssrpro=include
	set firewall.ssrpro.type=script
	set firewall.ssrpro.path=/etc/ssrpro.include
	set firewall.ssrpro.reload=1
	commit firewall
EOF

/etc/init.d/ssrpro stop
/etc/init.d/ssrpro enable

rm -f /tmp/luci-indexcache
exit 0
