# target/linux/ar71xx/base-files/etc/hotplug.d/net/10-ar922x-led-fix
#!/bin/sh

# For AR9220 and AR9223, GPIO JTAG must explicit be disabled
# before LEDs start working. Do this when wifi device is
# detected.

#
# $DEVPATH is not valid for some boards (including WZR-HP-AG300H).
# Manipulate the $DEVPATH to reach the corresponding phyN.
#

devdir=$(dirname $DEVPATH)
devdir=$(dirname $devdir)
phydir=/sys$devdir/ieee80211
phyname=$(cat $phydir/phy*/name)

[ -z $phyname -o $ACTION != "add" ] && exit 0

#
# ar922x_disable_gpio_jtag():
#
# Emulate
#       REG_SET_BIT(ah, AR_GPIO_INPUT_EN_VAL, AR_GPIO_JTAG_DISABLE);
# for AR9220 and AR9223.
#

ar922x_disable_gpio_jtag() {
	local regidx=0x4054

	[ -f /sys/kernel/debug/ieee80211/$1/ath9k/regidx ] || return

	echo $regidx > /sys/kernel/debug/ieee80211/$1/ath9k/regidx
	regval=$(cat /sys/kernel/debug/ieee80211/$1/ath9k/regval)
	regval=$((regval | 0x20000))
	echo regval $regval
	echo $regval > /sys/kernel/debug/ieee80211/$1/ath9k/regval
}

[ $phyname -a $ACTION = "add" ] && {
	. /lib/functions.sh

	case $(board_name) in
	wzr-hp-ag300h)
		ar922x_disable_gpio_jtag $phyname
		;;
	esac;
}

exit 0
