# target/linux/ramips/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

start() {
	. /lib/functions.sh

	local board=$(board_name)
	if [ $board = "sk-wb8" ]; then
		fw_setenv bootcount 0
	fi
}
