# package/lean/mt/drivers/mt7615d/files/etc/uci-defaults/10_mt7615_dbdc
#!/bin/sh

# fix isolate
sed -i 's/multicast_to_unicast:-1/multicast_to_unicast:-0/g' /lib/netifd/netifd-wireless.sh

# generate default wireless config
[ ! -f /etc/config/wireless ] && /sbin/wifi config

exit 0
