# package/lean/mt/drivers/mt7615d/files/etc/uci-defaults/mt7615-fix-isolate
#!/bin/sh

sed -i 's/multicast_to_unicast:-1/multicast_to_unicast:-0/g' /lib/netifd/netifd-wireless.sh
