# target/linux/omap24xx/base-files/etc/init.d/watchdog
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2010 OpenWrt.org

START=25

start() {
	[ -c /dev/watchdog -a -x /sbin/watchdog ] || {
		echo "WARNING: Watchdog not available. System will reboot soon!"
		return 1
	}
	/sbin/watchdog -T 63 -t 50 /dev/watchdog
	[ -x /usr/bin/schedtool ] && /usr/bin/schedtool -R -p 60 -n -20 $(pidof watchdog)
}

stop() {
	killall -q watchdog
}
