# customfeeds/packages/net/ola/files/olad.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009-2015 OpenWrt.org
# Copyright (C) 2015 Christian Beier <dontmind@freeshell.org>

USE_PROCD=1

START=90

start_service() {
	procd_open_instance
	procd_set_param user nobody
	procd_set_param command /usr/bin/olad --syslog --log-level 2 --config-dir /etc/ola
	procd_set_param respawn
	procd_close_instance
}
