# package/Twinzo/luci-app-mac/root/etc/uci-defaults/luci-mac
#!/bin/sh

uci -q batch <<-EOT
	set mac.@mac[0].interface="wan"
	set mac.@mac[0]."$(uci -q get network.wan_dev)"="wan_dev"
	set mac.@mac[0].macaddr="$(uci -q get network.wan.macaddr||uci -q get network.wan_dev.macaddr)"
	commit mac
    EOT

chmod 777 /etc/init.d/mac 

rm -f /tmp/luci-indexcache
exit 0
