# package/lean/luci-app-fullconenat/root/etc/uci-defaults/fullconenat
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@fullconenat[-1]
	add ucitrack fullconenat
	set ucitrack.@fullconenat[-1].init=fullconenat
	commit ucitrack
EOF

/etc/init.d/fullconenat enable

rm -f /tmp/luci-indexcache
exit 0
