# customfeeds/packages/net/tac_plus/files/tac_plus.init
#!/bin/sh /etc/rc.common

START=99

USE_PROCD=1
PROG=/usr/sbin/tac_plus

start_service() {
	procd_open_instance
	procd_set_param command "$PROG" -G -C /etc/tac_plus.conf
	procd_close_instance
}
