# customfeeds/packages/utils/lm-sensors/files/lm-sensors.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2017 Philip Prindeville, Redfish Solutions LLC

START=28
STOP=

PROG=/usr/sbin/sensors

start() {
	[ -f /etc/sensors.conf -o -f /etc/sensors3.conf ] || return

	$PROG -s
}

