# package/lean/luci-app-ngrokc/root/etc/uci-defaults/luci-app-ngrokc
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@ngrokc[-1]
	add ucitrack ngrokc
	set ucitrack.@ngrokc[-1].init=ngrokc
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
