# customfeeds/packages/libs/openldap/files/ldap.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009-2011 OpenWrt.org

START=60

SERVICE_USE_PID=1

start() {
	mkdir -m 0755 -p /var/openldap-data
	service_start /usr/sbin/slapd -h "ldap://localhost/ ldaps:///"
}

stop() {
	service_stop /usr/sbin/slapd
}

reload() {
	service_reload /usr/sbin/slapd
}
