# package/shadowvpn/files/shadowvpn.init
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2015 OpenWrt-dist
# Copyright (C) 2015 Jian Chang <aa65535@live.com>
#
# This is free software, licensed under the GNU General Public License v3.
# See /LICENSE for more information.
#

START=90
STOP=15

CONFIG=/var/etc/shadowvpn.conf

start_vpn() {
exit 0;
}

start() {
exit 0;
}

stop() {
exit 0;
}

boot() {
exit 0;
}
