# customfeeds/packages/utils/lcdproc/files/LCDd
#!/bin/sh /etc/rc.common

START=22
USE_PROCD=1
PROG=/usr/sbin/LCDd

config_file=/etc/LCDd.conf

start_service() {
	procd_open_instance
	procd_set_param command $PROG -f -c $config_file
	procd_close_instance
}

