# customfeeds/packages/utils/lcdproc/files/lcdproc
#!/bin/sh /etc/rc.common

START=50
USE_PROCD=1
PROG=/usr/bin/lcdproc

config_file=/etc/lcdproc.conf

start_service() {
	procd_open_instance
	procd_set_param command $PROG -f -c $config_file
	procd_close_instance
}

