# package/lean/luci-app-frpc/root/etc/uci-defaults/luci-frp
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@frp[-1]
	add ucitrack frp
	set ucitrack.@frp[-1].init=frp
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
