# target/linux/rockchip-rk3328/base-files/etc/init.d/fa-rk3328-led
#!/bin/sh /etc/rc.common

START=97
 
start() {
  # eth0
  echo netdev > /sys/class/leds/wan_led/trigger 
  echo eth0 > /sys/class/leds/wan_led/device_name 
  echo 1 > /sys/class/leds/wan_led/link
  # eth1
  echo netdev > /sys/class/leds/lan_led/trigger 
  echo eth1 > /sys/class/leds/lan_led/device_name 
  echo 1 > /sys/class/leds/lan_led/link
}

