# target/linux/rockchip-rk3328/base-files/etc/init.d/fa-rk3328-misc
#!/bin/sh /etc/rc.common

START=96
 
start() {
    [ -f /proc/irq/28/smp_affinity ] && echo 8 > /proc/irq/28/smp_affinity
    echo ondemand > /sys/devices/system/cpu/cpufreq/policy0/scaling_governor
}

