# customfeeds/luci/applications/luci-app-hd-idle/root/etc/uci-defaults/40_luci-hd_idle
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@hd-idle[-1]
	add ucitrack hd-idle
	set ucitrack.@hd-idle[-1].init=hd-idle
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
