# target/linux/rockchip-rk3399/base-files/etc/hotplug.d/iface/12-disable-rk3399-eth-offloading
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh

[ -n "$INTERFACE" ] || exit 3
[ "$INTERFACE" == "wan" ] || exit 4
[ "$ACTION" == "ifup" ] || exit 1

BOARD=$(cat /tmp/sysinfo/board_name | cut -d , -f2)
logger -t disable-offloading "board = $BOARD"
[ "$BOARD" == "nanopi-m4" -o "$BOARD" == "nanopc-t4" -o "$BOARD" == "nanopi-neo4" -o "$BOARD" == "som-rk3399" ] || exit 2
logger -t disable-offloading "rk3399 board detected"

ETHTOOL=/usr/sbin/ethtool
[ -f $ETHTOOL ] || exit 5

[ -d /sys/devices/platform/fe300000.ethernet/net/eth0 ] || exit 6
$ETHTOOL -K eth0 rx off tx off

logger -t disable-offloading "disabed rk3399 ethernet tcp/udp offloading tx/rx"
exit 0
