# target/linux/lantiq/falcon/base-files/etc/uci-defaults/01_led_migration
#
# Copyright (C) 2020 OpenWrt.org
#

. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
lantiq,easy98020)
	migrate_leds "easy98020:="
	;;
esac

remove_devicename_leds

migrations_apply system

exit 0
