# package/system/mountd/files/mountd.init
#!/bin/sh /etc/rc.common

START=80

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1
SERVICE_PID_FILE=/var/run/mountd.pid

MOUNTD_BIN=/sbin/mountd

start()
{
	P="$(uci get mountd.mountd.path)"
	[ -n "$P" -a ! -f "$P" ] && mkdir -p $P

	service_start $MOUNTD_BIN -f
}

stop()
{
	service_stop $MOUNTD_BIN
}
