# package/diy/3ginfo/files/etc/init.d/3ginfo
#!/bin/sh /etc/rc.common

# (c) 2010-2019 Cezary Jackiewicz <cezary@eko.one.pl>

START=99

start () {
	PORT=$(uci -q get 3ginfo.@3ginfo[0].http_port)
	[ -z $PORT ] && { PORT=81; }
	uhttpd -c /usr/share/3ginfo/uhttpd.conf -h /usr/share/3ginfo/ -p 0.0.0.0:$PORT > /dev/null 2>&1
}
