# package/community/openwrt-package/luci-app-control-webrestriction/root/etc/uci-defaults/luci-app-control-webrestriction
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@webrestriction[-1]
	add ucitrack webrestriction
	set ucitrack.@webrestriction[-1].init=webrestriction
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
