# package/lean/luci-app-unblockneteasemusic-mini/root/etc/uci-defaults/unblockneteasemusic-mini
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@unblockneteasemusic-mini[-1]
	add ucitrack unblockneteasemusic-mini
	set ucitrack.@unblockneteasemusic-mini[-1].init=unblockneteasemusic-mini
	commit ucitrack
	delete firewall.unblockneteasemusic_mini
	set firewall.unblockneteasemusic_mini=include
	set firewall.unblockneteasemusic_mini.type=script
	set firewall.unblockneteasemusic_mini.path=/var/etc/unblockneteasemusic-mini.include
	set firewall.unblockneteasemusic_mini.reload=1
	commit firewall
EOF

rm -f /tmp/luci-indexcache
exit 0
