# customfeeds/luci/applications/luci-app-dump1090/root/etc/uci-defaults/40_luci-dump1090
#!/bin/sh

# needed for "Save and Apply" to restart dump1090
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@dump1090[-1]
	add ucitrack dump1090
	set ucitrack.@dump1090[-1].init="dump1090"
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
