# target/linux/ipq40xx/base-files/etc/uci-defaults/04_led_migration
. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
engenius,emr3500)
	migrate_leds "emr3500:="
	;;
engenius,ens620ext|\
zyxel,nbg6617)
	migrate_leds ":wlan2G=:wlan2g" ":wlan5G=:wlan5g"
	;;
esac

remove_devicename_leds

migrations_apply system

exit 0
