# package/lean/luci-app-rclone/root/etc/uci-defaults/100_luci-rclone
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@rclone[-1]
	add ucitrack rclone
	set ucitrack.@rclone[-1].init=rclone
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
