# target/linux/ath79/generic/base-files/etc/uci-defaults/04_led_migration
#!/bin/sh

. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)
boardonly="${board##*,}"

case "$board" in
engenius,epg5000)
	migrate_leds ":wlan-2g=:wlan2g" ":wlan-5g=:wlan5g"
	;;
tplink,archer-c25-v1|\
tplink,archer-c58-v1|\
tplink,archer-c59-v1|\
tplink,archer-c59-v2|\
tplink,archer-c60-v1|\
tplink,archer-c60-v2|\
tplink,archer-c7-v4|\
tplink,archer-c7-v5|\
tplink,tl-wr902ac-v1)
	migrate_leds "^$boardonly:=tp-link:"
	;;
tplink,archer-c7-v2|\
tplink,tl-wdr3600-v1|\
tplink,tl-wdr4300-v1|\
tplink,tl-wdr4300-v1-il)
	migrate_leds ":blue:=:green:"
	;;
tplink,re355-v1)
	migrate_leds "re355:=tp-link:"
	;;
tplink,re450-v1)
	migrate_leds "re450:=tp-link:"
	;;
wd,mynet-n750)
	migrate_leds "wd:=mynet-n750:"
	;;
esac

migrations_apply system

exit 0
