# target/linux/ath79/generic/base-files/etc/uci-defaults/09_fix-checksum
. /lib/functions.sh

fix_seama_header() {
	local kernel_size=$(sed -n 's/mtd[0-9]*: \([0-9a-f]*\).*"kernel".*/\1/p' /proc/mtd)
	[ "$kernel_size" ] && mtd -c 0x$kernel_size fixseama firmware
}

fixwrgg() {
	local kernel_size=$(sed -n 's/mtd[0-9]*: \([0-9a-f]*\).*"kernel".*/\1/p' /proc/mtd)
	[ "$kernel_size" ] && mtd -c 0x$kernel_size fixwrgg firmware
}

board=$(board_name)

case "$board" in
dlink,dap-2695-a1)
	fixwrgg
	;;
qihoo,c301)
	fix_seama_header
	;;
esac
