# package/lean/luci-app-haproxy-tcp/root/etc/uci-defaults/99_luci-haproxy
#!/bin/sh

/etc/init.d/haproxy disable
/etc/init.d/haproxy stop

rm -f /tmp/luci-indexcache
exit 0
