# package/lean/mwan3-alt/files/etc/hotplug.d/iface/16-mwan3
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh
. /lib/mwan3/mwan3.sh

config_load mwan3
config_get_bool enabled globals 'enabled' '0'
[ ${enabled} -gt 0 ] || exit 0

if [ "$ACTION" == "ifup" ]; then
	mwan3_lock
	mwan3_rtmon
	mwan3_unlock
fi

exit 0
