# target/linux/rockchip/rk35xx/base-files/etc/init.d/reload-onboard-wifi
#!/bin/sh /etc/rc.common

START=99

boot() {
	case "$(board_name)" in
	armsom,sige1)
		( sleep 50s; wifi; ) &
		;;
	armsom,sige7-v1)
		( sleep 100s; wifi; ) &
		;;
	esac
}
