# target/linux/ipq807x/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	redmi,ax6|\
	xiaomi,ax3600)
		# OTA handling should not be used. Reset it just in case.
		fw_setenv flag_ota_reboot 0
		# Not strictly needed but useful to handle partition crash condition
		fw_setenv flag_boot_success 1
	;;
	esac
	
	for q in $(ls /sys/class/net/eth*/queues/rx-*/rps_cpus); do echo f > $q; done
	for q in $(ls /sys/class/net/eth*/queues/rx-*/rps_flow_cnt); do echo 4096 > $q; done
	for q in $(ls /sys/class/net/eth*/queues/tx-*/xps_cpus); do echo f > $q; done
	echo 32768 > /proc/sys/net/core/rps_sock_flow_entries
}
