# target/linux/ipq807x/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	redmi,ax6|\
	xiaomi,ax3600)
		# OTA handling should not be used. Reset it just in case.
		fw_setenv flag_ota_reboot 0
		# Not strictly needed but useful to handle partition crash condition
		fw_setenv flag_boot_success 1
	;;
	esac
}
