# customfeeds/packages/admin/gkrellmd/files/gkrellmd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org

START=60
BIN=gkrellmd
RUN_D=/var/run
PID_F=$RUN_D/$BIN.pid

start() {
	mkdir -p $RUN_D
	$BIN $OPTIONS
}

stop() {
	[ -f $PID_F ] && kill $(cat $PID_F)
}
