# package/Twinzo/DevicesComment/root/etc/uci-defaults/luci-devcom
#!/bin/sh

chmod 755 /etc/init.d/devcom 

rm -f /tmp/luci-indexcache
exit 0