# package/lean/luci-app-easymesh/root/etc/uci-defaults/luci-easymesh
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@easymesh[-1]
	add ucitrack easymesh
	set ucitrack.@easymesh[-1].init=easymesh
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
