# target/linux/mediatek/base-files/etc/init.d/mtkhnat
#!/bin/sh /etc/rc.common

START=90

USE_PROCD=1
NAME=mtkhnat
PROG=/sbin/mtkhnat

start_service() {
	procd_open_instance
	procd_set_param command "${PROG}"
	procd_close_instance
}
