# package/diy/luci-app-cpe/root/etc/uci-defaults/luci-modem
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@modem[-1]
	add ucitrack modem
	set ucitrack.@modem[-1].init=modem
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
