# package/my/luci-app-wrtbwmon/root/etc/hotplug.d/iface/60-wrtbwmon
#!/bin/sh
# Copyright (C) 2007 OpenWrt.org

[ "$ACTION" = "ifup" ] || exit 0
/etc/init.d/wrtbwmon start
