# target/linux/lantiq/ase/base-files/etc/uci-defaults/01_led_migration
#
# Copyright (C) 2020 OpenWrt.org
#

. /lib/functions/migrations.sh

remove_devicename_leds

migrations_apply system

exit 0
