# target/linux/apm821xx/base-files/etc/uci-defaults/04_led_migration
. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
meraki,mr24)
	migrate_leds "orange:power=amber:fault" \
		":wifi1=:wlan-0" \
		":wifi2=:wlan-1" \
		":wifi3=:wlan-2" \
		":wifi4=:wlan-3"
	;;
meraki,mx60)
	migrate_leds ":lan1=:lan-0" \
		":lan2=:lan-1" \
		":lan3=:lan-2" \
		":lan4=:lan-3" \
		"orange:power=amber:power" \
		"orange:wan=amber:wan"
	;;
netgear,wndap620)
	migrate_leds ":activity=:lan-0" \
		":test=:fault" \
		":wlan2g=:wlan-0" \
		":wlan5g=:wlan-1" \
		":link100=:lan" \
		":link1000=:lan-1"
	;;
netgear,wndap660)
	migrate_leds ":activity=:lan-0" \
		":test=:fault" \
		":wlan2g=:wlan-0" \
		":wlan5g=:wlan-1" \
		":lan1-link100=:lan-0" \
		":lan1-link1000=:lan-1" \
		":lan2-link100=:lan-1" \
		":lan2-link1000=:lan-2"
	;;
netgear,wndr4700)
	migrate_leds "orange:power=amber:fault" \
		"white:logo=white:indicator" \
		"green:hd=green:disk" \
		"red:hd=red:disk-err"
	;;
wd,mybooklive)
	migrate_leds "red:power=red:fault" \
		"blue:power=blue:disk"
	;;
esac

remove_devicename_leds

migrations_apply system

exit 0
