# package/lean/luci-app-verysync/root/etc/uci-defaults/luci-verysync
#!/bin/sh
touch /etc/config/verysync

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@verysync[-1]
	add ucitrack verysync
	set ucitrack.@verysync[-1].exec='/etc/init.d/verysync restart'
	commit ucitrack
EOF
# remove LuCI cache
rm -f /tmp/luci*
exit 0
