# package/lean/luci-app-ps3netsrv/root/etc/uci-defaults/ps3netsrv
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@ps3netsrv[-1]
	add ucitrack ps3netsrv
	set ucitrack.@ps3netsrv[-1].init=ps3netsrv
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
