# customfeeds/packages/utils/bluez/files/bluetoothd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org

#start after dbus (60)
START=62
USE_PROCD=1
PROG=/usr/bin/bluetoothd

start_service() {
	ln -snf /etc/bluetooth/keys/ /var/lib/bluetooth
	procd_open_instance
	procd_set_param command "$PROG" -n
	procd_close_instance
}
