# target/linux/ipq40xx/base-files/etc/hotplug.d/firmware/11-ath10k-caldata
#!/bin/sh

# xor multiple hex values of the same length
xor() {
	local val
	local ret="0x$1"
	local retlen=${#1}

	shift
	while [ -n "$1" ]; do
		val="0x$1"
		ret=$((ret ^ val))
		shift
	done

	printf "%0${retlen}x" "$ret"
}

ath10kcal_die() {
	echo "ath10cal: " "$*"
	exit 1
}

ath10kcal_extract() {
	local part=$1
	local offset=$2
	local count=$3
	local mtd

	mtd=$(find_mtd_chardev $part)
	[ -n "$mtd" ] || \
		ath10kcal_die "no mtd device found for partition $part"

	dd if=$mtd of=/lib/firmware/$FIRMWARE iflag=skip_bytes bs=$count skip=$offset count=1 2>/dev/null || \
		ath10kcal_die "failed to extract calibration data from $mtd"
}

ath10kcal_ubi_extract() {
	local part=$1
	local offset=$2
	local count=$3
	local ubidev
	local ubi

	. /lib/upgrade/nand.sh

	ubidev=$(nand_find_ubi $CI_UBIPART)
	ubi=$(nand_find_volume $ubidev $part)
	[ -n "$ubi" ] || \
		ath10kcal_die "no UBI volume found for $part"

	dd if=/dev/$ubi of=/lib/firmware/$FIRMWARE iflag=skip_bytes bs=$count skip=$offset count=1 2>/dev/null || \
		ath10kcal_die "failed to extract from $ubi"
}

ath10kcal_patch_mac_crc() {
	local mac=$1
	local mac_offset=6
	local chksum_offset=2
	local xor_mac
	local xor_fw_mac
	local xor_fw_chksum

	[ -z "$mac" ] && return

	xor_fw_mac=$(hexdump -v -n 6 -s $mac_offset -e '/1 "%02x"' /lib/firmware/$FIRMWARE)
	xor_fw_mac="${xor_fw_mac:0:4} ${xor_fw_mac:4:4} ${xor_fw_mac:8:4}"

	macaddr_2bin $mac | dd of=/lib/firmware/$FIRMWARE conv=notrunc bs=1 seek=6 count=6

	xor_mac=${mac//:/}
	xor_mac="${xor_mac:0:4} ${xor_mac:4:4} ${xor_mac:8:4}"

	xor_fw_chksum=$(hexdump -v -n 2 -s $chksum_offset -e '/1 "%02x"' /lib/firmware/$FIRMWARE)
	xor_fw_chksum=$(xor $xor_fw_chksum $xor_fw_mac $xor_mac)

	printf "%b" "\x${xor_fw_chksum:0:2}\x${xor_fw_chksum:2:2}" | \
		dd of=/lib/firmware/$FIRMWARE conv=notrunc bs=1 seek=$chksum_offset count=2
}

ath10kcal_is_caldata_valid() {
	local expected="$1"

	magic=$(hexdump -v -n 2 -e '1/1 "%02x"' /lib/firmware/$FIRMWARE)
	[[ "$magic" == "$expected" ]]
	return $?
}

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)


case "$FIRMWARE" in
"ath10k/cal-pci-0000:01:00.0.bin")
	case "$board" in
	meraki,mr33)
		ath10kcal_ubi_extract "ART" 36864 2116
		ath10kcal_is_caldata_valid "4408" || ath10kcal_extract "ART" 36864 2116
		ath10kcal_patch_mac_crc $(macaddr_add $(get_mac_binary "/sys/bus/i2c/devices/0-0050/eeprom" 102) +1)
		;;
	esac
	;;
"ath10k/pre-cal-pci-0000:01:00.0.bin")
	case "$board" in
	asus,map-ac2200)
		ath10kcal_ubi_extract "Factory" 36864 12064
		ln -sf /lib/firmware/ath10k/pre-cal-pci-0000\:00\:00.0.bin \
			/lib/firmware/ath10k/QCA9888/hw2.0/board.bin
		;;
	asus,rt-acrh17)
		CI_UBIPART=UBI_DEV
		ath10kcal_ubi_extract "Factory" 36864 12064
		;;
	openmesh,a62)
		ath10kcal_extract "0:ART" 36864 12064
		;;
	esac
	;;
"ath10k/pre-cal-ahb-a000000.wifi.bin")
	case "$board" in
	8dev,jalapeno |\
	glinet,gl-b1300 |\
	linksys,ea6350v3 |\
	qcom,ap-dk01.1-c1)
		ath10kcal_extract "ART" 4096 12064
		;;
	asus,map-ac2200)
		ath10kcal_ubi_extract "Factory" 4096 12064
		;;
	p2w,r619ac |\
	asus,rt-acrh17|\
	asus,rt-ac58u)
		CI_UBIPART=UBI_DEV
		ath10kcal_ubi_extract "Factory" 4096 12064
		;;
	avm,fritzbox-4040)
		/usr/bin/fritz_cal_extract -i 1 -s 0x400 -e 0x207 -l 12064 -o /lib/firmware/$FIRMWARE $(find_mtd_chardev "urlader_config")
		;;
	meraki,mr33)
		ath10kcal_ubi_extract "ART" 4096 12064
		ath10kcal_is_caldata_valid "202f" || ath10kcal_extract "ART" 4096 12064
		ath10kcal_patch_mac_crc $(macaddr_add $(get_mac_binary "/sys/bus/i2c/devices/0-0050/eeprom" 102) +2)
		;;
	netgear,ex6100v2 |\
	netgear,ex6150v2)
		ath10kcal_extract "ART" 4096 12064
		ath10kcal_patch_mac_crc $(mtd_get_mac_binary dnidata 0)
		;;
	compex,wpj428 |\
	engenius,eap1300 |\
	hugo,ac1200 |\
	openmesh,a42 |\
	openmesh,a62)
		ath10kcal_extract "0:ART" 4096 12064
		;;
	zyxel,nbg6617 |\
	zyxel,wre6606)
		ath10kcal_extract "ART" 4096 12064
		ath10kcal_patch_mac_crc $(macaddr_add $(cat /sys/class/net/eth0/address) -2)
		;;
	esac
	;;
"ath10k/pre-cal-ahb-a800000.wifi.bin")
	case "$board" in
	8dev,jalapeno |\
	glinet,gl-b1300 |\
	linksys,ea6350v3 |\
	p2w,r619ac |\
	qcom,ap-dk01.1-c1)
		ath10kcal_extract "ART" 20480 12064
		;;
	asus,map-ac2200)
		ath10kcal_ubi_extract "Factory" 20480 12064
		;;
	asus,rt-ac58u)
		CI_UBIPART=UBI_DEV
		ath10kcal_ubi_extract "Factory" 20480 12064
		;;
	avm,fritzbox-4040)
		/usr/bin/fritz_cal_extract -i 1 -s 0x400 -e 0x208 -l 12064 -o /lib/firmware/$FIRMWARE $(find_mtd_chardev "urlader_config")
		;;
	meraki,mr33)
		ath10kcal_ubi_extract "ART" 20480 12064
		ath10kcal_is_caldata_valid "202f" || ath10kcal_extract "ART" 20480 12064
		ath10kcal_patch_mac_crc $(macaddr_add $(get_mac_binary "/sys/bus/i2c/devices/0-0050/eeprom" 102) +3)
		;;
	netgear,ex6100v2 |\
	netgear,ex6150v2)
		ath10kcal_extract "ART" 20480 12064
		ath10kcal_patch_mac_crc $(mtd_get_mac_binary dnidata 12)
		;;
	compex,wpj428 |\
	engenius,eap1300 |\
	hugo,ac1200 |\
	openmesh,a42 |\
	openmesh,a62)
		ath10kcal_extract "0:ART" 20480 12064
		;;
	zyxel,nbg6617 |\
	zyxel,wre6606)
		ath10kcal_extract "ART" 20480 12064
		ath10kcal_patch_mac_crc $(macaddr_add $(cat /sys/class/net/eth0/address) -1)
		;;
	esac
	;;
*)
	exit 1
	;;
esac
