# target/linux/ipq40xx/base-files/etc/hotplug.d/firmware/11-ath10k-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"ath10k/cal-pci-0000:01:00.0.bin")
	case "$board" in
	meraki,mr33)
		caldata_extract_ubi "ART" 0x9000 0x844
		caldata_valid "4408" || caldata_extract "ART" 0x9000 0x844
		ath10k_patch_mac $(macaddr_add $(get_mac_binary "/sys/bus/i2c/devices/0-0050/eeprom" 0x66) +1)
		;;
	esac
	;;
"ath10k/pre-cal-pci-0000:01:00.0.bin")
	case "$board" in
	asus,map-ac2200)
		caldata_extract_ubi "Factory" 0x9000 0x2f20
		ln -sf /lib/firmware/ath10k/pre-cal-pci-0000\:00\:00.0.bin \
			/lib/firmware/ath10k/QCA9888/hw2.0/board.bin
		;;
	asus,rt-acrh17)
		CI_UBIPART=UBI_DEV
		caldata_extract_ubi "Factory" 0x9000 0x2f20
		;;
	openmesh,a62)
		caldata_extract "0:ART" 0x9000 0x2f20
		;;
	p2w,r619ac |\
	p2w,r619ac-128m)
		caldata_extract "ART" 0x9000 0x2f20
		;;
	esac
	;;
"ath10k/pre-cal-ahb-a000000.wifi.bin")
	case "$board" in
	8dev,jalapeno |\
	glinet,gl-b1300 |\
	linksys,ea6350v3 |\
	p2w,r619ac |\
	p2w,r619ac-128m |\
	qcom,ap-dk01.1-c1)
		caldata_extract "ART" 0x1000 0x2f20
		;;
	asus,map-ac2200)
		caldata_extract_ubi "Factory" 0x1000 0x2f20
		;;
	asus,rt-acrh17|\
	asus,rt-ac58u)
		CI_UBIPART=UBI_DEV
		caldata_extract_ubi "Factory" 0x1000 0x2f20
		;;
	avm,fritzbox-4040)
		/usr/bin/fritz_cal_extract -i 1 -s 0x400 -e 0x207 -l 12064 -o /lib/firmware/$FIRMWARE $(find_mtd_chardev "urlader_config")
		;;
	compex,wpj428 |\
	engenius,eap1300 |\
	hugo,ac1200 |\
	openmesh,a42 |\
	openmesh,a62)
		caldata_extract "0:ART" 0x1000 0x2f20
		;;
	meraki,mr33)
		caldata_extract_ubi "ART" 0x1000 0x2f20
		caldata_valid "202f" || caldata_extract 0x1000 0x2f20
		ath10k_patch_mac $(macaddr_add $(get_mac_binary "/sys/bus/i2c/devices/0-0050/eeprom" 0x66) +2)
		;;
	netgear,ex6100v2 |\
	netgear,ex6150v2)
		caldata_extract "ART" 0x1000 0x2f20
		ath10k_patch_mac $(mtd_get_mac_binary dnidata 0x0)
		;;
	zyxel,nbg6617 |\
	zyxel,wre6606)
		caldata_extract "ART" 0x1000 0x2f20
		ath10k_patch_mac $(macaddr_add $(cat /sys/class/net/eth0/address) -2)
		;;
	esac
	;;
"ath10k/pre-cal-ahb-a800000.wifi.bin")
	case "$board" in
	8dev,jalapeno |\
	glinet,gl-b1300 |\
	linksys,ea6350v3 |\
	p2w,r619ac |\
	p2w,r619ac-128m |\
	qcom,ap-dk01.1-c1)
		caldata_extract "ART" 0x5000 0x2f20
		;;
	asus,map-ac2200)
		caldata_extract "Factory" 0x5000 0x2f20
		;;
	asus,rt-ac58u)
		CI_UBIPART=UBI_DEV
		caldata_extract_ubi "Factory" 0x5000 0x2f20
		;;
	avm,fritzbox-4040)
		/usr/bin/fritz_cal_extract -i 1 -s 0x400 -e 0x208 -l 12064 -o /lib/firmware/$FIRMWARE $(find_mtd_chardev "urlader_config")
		;;
	compex,wpj428 |\
	engenius,eap1300 |\
	hugo,ac1200 |\
	openmesh,a42 |\
	openmesh,a62)
		caldata_extract "0:ART" 0x5000 0x2f20
		;;
	meraki,mr33)
		caldata_extract_ubi "ART" 0x5000 0x2f20
		caldata_valid "202f" || caldata_extract "ART" 0x5000 0x2f20
		ath10k_patch_mac $(macaddr_add $(get_mac_binary "/sys/bus/i2c/devices/0-0050/eeprom" 0x66) +3)
		;;
	netgear,ex6100v2 |\
	netgear,ex6150v2)
		caldata_extract "ART" 0x5000 0x2f20
		ath10k_patch_mac $(mtd_get_mac_binary dnidata 0xc)
		;;
	zyxel,nbg6617 |\
	zyxel,wre6606)
		caldata_extract "ART" 0x5000 0x2f20
		ath10k_patch_mac $(macaddr_add $(cat /sys/class/net/eth0/address) -1)
		;;
	esac
	;;
*)
	exit 1
	;;
esac
