# customfeeds/packages/net/apache/files/apache2.init
#!/bin/sh /etc/rc.common

START=99

reload() {
	apachectl -k restart
}

restart() {
	apachectl -k restart
}

start() {
	mkdir -p -m 0750 /var/log/apache2
	mkdir -p /var/run/apache2
	apachectl -k start
}

stop() {
	apachectl -k stop
}
