# customfeeds/luci/applications/luci-app-aria2/root/etc/uci-defaults/40_luci-aria2
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@aria2[-1]
	add ucitrack aria2
	set ucitrack.@aria2[-1].init=aria2
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
