# package/lean/luci-app-ipsec-vpnd/root/etc/uci-defaults/ipsec
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete firewall.ipsecd
	set firewall.ipsecd=include
	set firewall.ipsecd.type=script
	set firewall.ipsecd.path=/etc/ipsec.include
	set firewall.ipsecd.reload=1
	commit firewall
EOF

/etc/init.d/ipsec disable && /etc/init.d/ipsec stop

rm -f /tmp/luci-indexcache
exit 0
