# target/linux/ipq40xx/base-files/etc/hotplug.d/net/21_adjust_network
#!/bin/sh

[ -f /lib/adjust_network.sh ] && {
	. /lib/adjust_network.sh

	adjust_eth_queue
}
