# package/lean/luci-app-jd-dailybonus/root/etc/init.d/jd-dailybonus
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2020 Jerryk
#
# This is free software, licensed under the GNU General Public License v3.
# See /LICENSE for more information.
#

USE_PROCD=1
START=25

NAME=jd-dailybonus
CRON_FILE=/etc/crontabs/root
JD_SCRIPT=/usr/share/jd-dailybonus/JD_DailyBonus.js

uci_get_by_name() {
	local ret=$(uci get $NAME.$1.$2 2>/dev/null)
	echo ${ret:=$3}
}

uci_get_by_type() {
	local ret=$(uci get $NAME.@$1[0].$2 2>/dev/null)
	echo ${ret:=$3}
}


add_cron() {
	sed -i '/jd-dailybonus/d' $CRON_FILE
	[ $(uci_get_by_type global auto_run 0) -eq 1 ] && echo '5 '$(uci_get_by_type global auto_run_time)' * * * /bin/bash -c "sleep $[RANDOM % 180]s"; /usr/share/jd-dailybonus/newapp.sh -r' >> $CRON_FILE
	[ $(uci_get_by_type global auto_update 0) -eq 1 ] && echo '1 '$(uci_get_by_type global auto_update_time)' * * * /usr/share/jd-dailybonus/newapp.sh -u' >> $CRON_FILE
	crontab $CRON_FILE
}

del_cron() {
	sed -i '/jd-dailybonus/d' $CRON_FILE
	/etc/init.d/cron restart
}

fill_cookie() {
    cookie1=$(uci_get_by_type global cookie)
    if [ ! "$cookie1" = "" ]; then
        varb="var Key = '$cookie1'" 
        sed -i "s/^var Key =.*/$varb/g" $JD_SCRIPT
    fi

    cookie2=$(uci_get_by_type global cookie2)
    if [ ! "$cookie2" = "" ]; then
        varb2="var DualKey = '$cookie2'" 
        sed -i "s/^var DualKey =.*/$varb2/g" $JD_SCRIPT
    fi
}

start_service(){
	fill_cookie
	add_cron
}

stop_service() {
    del_cron
}

service_triggers() {
	procd_add_reload_trigger "jd-dailybonus"
}