# package/lean/shellinabox/files/shellinaboxd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 snqu.com

USE_PROCD=1
START=99

BIN=/usr/sbin/shellinaboxd

start_service() {
	procd_open_instance
	procd_set_param command $BIN -u root -c /etc/shellinabox \
		--css=/etc/shellinabox/black-on-white.css \
		--css=/etc/shellinabox/white-on-black.css \
		--user-css Normal:+/etc/shellinabox/black-on-white.css,Reverse:-/etc/shellinabox/white-on-black.css
	procd_close_instance
}