# package/diy/3ginfo/files-text/etc/hotplug.d/usb/99-3ginfo
case "$ACTION" in
	remove)
		[ -f /tmp/pincode_was_given ] && rm -f /tmp/pincode_was_given;;
esac
