# target/linux/ath79/mikrotik/base-files/etc/hotplug.d/firmware/10-ath9k-eeprom
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

wlan_data="/sys/firmware/mikrotik/hard_config/wlan_data"
mac_base="$(cat /sys/firmware/mikrotik/hard_config/mac_base)"

board=$(board_name)

case "$FIRMWARE" in
"ath9k-eeprom-ahb-18100000.wmac.bin")
	case $board in
	mikrotik,routerboard-wap-g-5hact2hnd)
		caldata_from_file $wlan_data 0x1000 0x440 /tmp/$FIRMWARE
		ath9k_patch_mac $(macaddr_add "$mac_base" +2) /tmp/$FIRMWARE
		caldata_sysfsload_from_file /tmp/$FIRMWARE 0x0 0x440
		rm -f /tmp/$FIRMWARE
		;;
	*)
		caldata_die "board $board is not supported yet"
		;;
	esac
	;;
esac
