# customfeeds/luci/applications/luci-app-adblock/root/etc/uci-defaults/40_luci-adblock
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@adblock[-1]
	add ucitrack adblock
	set ucitrack.@adblock[-1].init=adblock
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
