# package/lean/luci-app-netspeedtest/root/etc/uci-defaults/luci-netspeedtest
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@netspeedtest[-1]
	add ucitrack netspeedtest
	set ucitrack.@netspeedtest[-1].init=netspeedtest
	commit ucitrack
EOF

rm -rf /tmp/luci-*
exit 0
