# customfeeds/luci/applications/luci-app-upnp/root/etc/uci-defaults/40_luci-miniupnp
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@upnpd[-1]
	add ucitrack upnpd
	set ucitrack.@upnpd[-1].init=miniupnpd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
