# customfeeds/packages/net/beanstalkd/files/beanstalkd.init
#!/bin/sh /etc/rc.common
# Beanstalk
 
START=21
STOP=11

USE_PROCD=1
PROG=/usr/bin/beanstalkd
 
start_service() {
    procd_open_instance
    procd_set_param command $PROG
    procd_close_instance
}
