# package/lean/autocore/files/loongarch64/autocore
#!/bin/sh /etc/rc.common
# Copyright (C) 2025 lean <coolsnowwolf@gmail.com>

START=99

start()
{
	a=$(cat /proc/cpuinfo | grep 'Model Name' | cut -f2 -d: | head -n 1)
	b=$(echo -n ' : ')
	c=$(cat /proc/cpuinfo | grep 'core' | sort -u | wc -l)
	d=$(echo -n 'C')
	e=$(cat /proc/cpuinfo | grep 'global_id' | sort -u | wc -l)
	f=$(echo -n 'T ')
	g=$(cat /proc/cpuinfo | grep 'CPU Family' | cut -f2 -d: | head -n 1)
	
	h=${g}' - '${a}${b}${c}${d}${e}${f}

	mkdir -p /tmp/sysinfo
	echo $h > /tmp/sysinfo/model
}
