# package/lean/softethervpn/files/vpnserver.init
#!/bin/sh /etc/rc.common
START=60

start(){
        logger -t 'softethervpn' "Starting softether vpnserver service."
	/usr/bin/env LANG=en_US.UTF-8 /usr/libexec/softethervpn/vpnserver start
}

stop(){
        logger -t 'softethervpn' "Stopping softether vpnserver service."
	/usr/bin/env LANG=en_US.UTF-8 /usr/libexec/softethervpn/vpnserver stop
}
