# package/community/openwrt-package/luci-app-syncthing/root/etc/uci-defaults/luci-syncthing
#!/bin/sh
touch /etc/config/syncthing

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@syncthing[-1]
	add ucitrack syncthing
	set ucitrack.@syncthing[-1].exec='/etc/init.d/syncthing restart'
	commit ucitrack
EOF
# remove LuCI cache
rm -f /tmp/luci*
exit 0
