# target/linux/gemini/base-files/etc/uci-defaults/09_fix-checksum
#
# Copyright (C) 2019 OpenWrt.org
#

. /lib/functions.sh

board=$(board_name)

fixwrgg() {
	local kernel_size=$(sed -n 's/mtd[0-9]*: \([0-9a-f]*\).*"kernel".*/\1/p' /proc/mtd)

	[ "$kernel_size" ] && mtd -c 0x$kernel_size fixwrgg firmware
}

case "$board" in
dlink,dir-685)
	fixwrgg
	;;
esac
