# customfeeds/packages/net/ostiary/files/ostiaryd.init
#!/bin/sh /etc/rc.common

START=99
STOP=85

USE_PROCD=1

start_service()
{

	procd_open_instance
	procd_set_param command /usr/bin/ostiaryd -c /etc/ostiary.cfg 
	procd_close_instance
}


reload_service() {
         procd_send_signal ostiaryd
}
