# package/lean/luci-app-dnscrypt-proxy/files/root/etc/uci-defaults/dnscrypt-proxy
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@dnscrypt-proxy[-1]
	add ucitrack dnscrypt-proxy
	set ucitrack.@dnscrypt-proxy[-1].init=dnscrypt-proxy
	commit ucitrack
	delete firewall.dnscrypt-proxy
	set firewall.dnscrypt-proxy=include
	set firewall.dnscrypt-proxy.type=script
	set firewall.dnscrypt-proxy.path=/var/etc/dnscrypt-proxy.include
	set firewall.dnscrypt-proxy.reload=1
	commit firewall
EOF

rm -f /tmp/luci-indexcache
exit 0
