# package/lean/luci-app-mosdns/root/etc/hotplug.d/iface/99-mosdns
#!/bin/sh
[ "$ACTION" = ifup ] && /etc/init.d/mosdns restart
