# target/linux/ipq806x/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

. "$IPKG_INSTROOT/lib/upgrade/asrock.sh"

boot() {
	case $(board_name) in
	asrock,g10)
		asrock_bootconfig_mangle "bootcheck" && reboot
		;;
	edgecore,ecw5410)
		fw_setenv bootcount 0
		;;
	linksys,ea7500-v1 |\
	linksys,ea8500)
		mtd resetbc s_env || true
		;;
	esac
	
	echo 800000 > /sys/devices/system/cpu/cpu0/cpufreq/scaling_min_freq
	echo 800000 > /sys/devices/system/cpu/cpu1/cpufreq/scaling_min_freq
	echo 20 > /sys/devices/system/cpu/cpufreq/ondemand/up_threshold
	echo 10 > /sys/devices/system/cpu/cpufreq/ondemand/sampling_down_factor
	
	sysctl -w vm.min_free_kbytes=65536
	sysctl -w net.netfilter.nf_conntrack_max=32768
}
