# target/linux/omap/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyO0::askfirst:/usr/libexec/login.sh
ttyO2::askfirst:/usr/libexec/login.sh
tty1::askfirst:/usr/libexec/login.sh
