# customfeeds/packages/net/strongswan/files/etc/hotplug.d/ipsec/01-user
#!/bin/sh

[ -e "/etc/ipsec.user" ] && {
	. /etc/ipsec.user "$2"
}

exit 0
