# package/lean/luci-app-dogcom/files/root/etc/uci-defaults/40_luci-dogcom
#!/bin/sh

uci -q batch <<-EOF >/dev/null
    delete ucitrack.@dogcom[-1]
    add ucitrack dogcom
    set ucitrack.@dogcom[-1].init=dogcom
    commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0