# customfeeds/packages/net/dynapoint/src/dynapoint.init
#!/bin/sh /etc/rc.common

START=99

USE_PROCD=1
PROG=/usr/sbin/dynapoint.lua
CONFFILE=/etc/config/dynapoint

start_service() {
	procd_open_instance
	procd_set_param command $PROG
	procd_set_param file $CONFFILE
	procd_set_param respawn
	procd_close_instance
}
