# package/Twinzo/luci-app-devcom/root/etc/uci-defaults/luci-devcom
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@devcom[-1]
	add ucitrack devcom
	set ucitrack.@devcom[-1].init=devcom
	commit ucitrack
EOF

chmod 755 /etc/init.d/devcom 

rm -f /tmp/luci-indexcache
exit 0