# target/linux/ramips/base-files/etc/hotplug.d/firmware/10-rt2x00-eeprom
#!/bin/sh

rt2x00_eeprom_die() {
	echo "rt2x00 eeprom: " "$*"
	exit 1
}

rt2x00_eeprom_extract() {
	local part=$1
	local offset=$2
	local count=$3
	local mtd

	mtd=$(find_mtd_part $part)
	[ -n "$mtd" ] || \
		rt2x00_eeprom_die "no mtd device found for partition $part"

	dd if=$mtd of=/lib/firmware/$FIRMWARE bs=1 skip=$offset count=$count 2>/dev/null || \
		rt2x00_eeprom_die "failed to extract from $mtd"
}

jboot_eeprom_extract() {
	local part=$1
	local offset=$2
	local mtd

	mtd=$(find_mtd_part $part)
	[ -n "$mtd" ] || \
		rt2x00_eeprom_die "no mtd device found for partition $part"

	jboot_config_read -i $mtd -o $offset -e /lib/firmware/$FIRMWARE  2>/dev/null || \
		rt2x00_eeprom_die "failed to extract from $mtd"
}

rt2x00_eeprom_set_macaddr() {
	local macaddr=$1

	[ -n "$macaddr" ] || \
		rt2x00_eeprom_die "invalid wlan mac address"

	macaddr_2bin $macaddr | dd of=/lib/firmware/$FIRMWARE \
				conv=notrunc bs=1 seek=4 count=6 2>/dev/null || \
		rt2x00_eeprom_die "failed to write mac address to eeprom file"
}

FW="/lib/firmware/$FIRMWARE"
[ -e "$FW" ] && exit 0

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)

case "$FIRMWARE" in
"soc_wmac.eeprom")
	case $board in
	dlink,dwr-116-a1|\
	dlink,dwr-118-a2|\
	dlink,dwr-921-c1|\
	lava,lr-25g001)
		wan_mac=$(jboot_config_read -m -i $(find_mtd_part "config") -o 0xE000)
		wifi_mac=$(macaddr_add "$wan_mac" 1)
		jboot_eeprom_extract "config" 0xE000
		rt2x00_eeprom_set_macaddr $wifi_mac
		;;
	tiny-ac)
		wifi_mac=$(mtd_get_mac_ascii u-boot-env INIC_MAC_ADDR)
		rt2x00_eeprom_extract "factory" 0 512
		rt2x00_eeprom_set_macaddr $wifi_mac
		;;
	*)
		rt2x00_eeprom_die "Please define mtd-eeprom in $board DTS file!"
		;;
	esac
	;;
esac
