# target/linux/mpc85xx/base-files/etc/uci-defaults/05_fix-compat-version
. /lib/functions.sh
. /lib/functions/system.sh

case "$(board_name)" in
aerohive,hiveap-330)
	uci set system.@system[0].compat_version="2.0"
	uci commit system
	;;
esac

exit 0
