# package/lean/luci-app-dnspod/files/root/etc/init.d/dnspod
#!/bin/sh /etc/rc.common
START=80

run_reboot()
{                                                                                                       
	killall dnspod.sh
	/sbin/dnspod.sh --svc
}

start()
{
	/sbin/dnspod.sh --svc
}

stop()
{
	killall dnspod.sh
}
