# package/lean/luci-app-dnspod/files/root/etc/init.d/dnspod
#!/bin/sh /etc/rc.common
START=80

start()
{
	/usr/sbin/dnspod.sh --svc &
}

stop()
{
	ps | grep dnspod.sh | grep -v 'grep' | awk '{print $1}' | xargs kill
}

run_reboot()
{
	stop
	start
}
