# customfeeds/luci/applications/luci-app-nlbwmon/root/etc/uci-defaults/40_luci-nlbwmon
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@nlbwmon[-1]
	add ucitrack nlbwmon
	set ucitrack.@nlbwmon[-1].init=nlbwmon
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
