# target/linux/ramips/rt305x/base-files/etc/uci-defaults/04_led_migration
. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
fon,fonera-20n)
	migrate_leds "^fonera20n:=fonera-20n:"
	;;
intenso,memory2move)
	migrate_leds "^m2m:=memory2move:"
	;;
esac

migrations_apply system

exit 0
