# package/lean/shadowsocksR-libev-full/files/shadowsocksr
#!/bin/sh /etc/rc.common

START=95

SERVICE_USE_PID=1
SERVICE_WRITE_PID=1
SERVICE_DAEMONIZE=1

CONFIG=/etc/shadowsocksr.json

start() {
	service_start /usr/bin/ssr-local -c $CONFIG -l 8888
	#service_start /usr/bin/ssr-redir -c $CONFIG
	#service_start /usr/bin/ssr-tunnel -c $CONFIG -l 5353 -L 8.8.8.8:53 -U
}

stop() {
	service_stop /usr/bin/ssr-local
	#service_stop /usr/bin/ssr-redir
	#service_stop /usr/bin/ssr-tunnel
}
