# target/linux/bcm63xx/base-files/etc/uci-defaults/04_led_migration
. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
brcm,bcm96318ref-p300|\
brcm,bcm963281tan|\
brcm,bcm96328avng|\
d-link,dsl-2640b-b|\
d-link,dva-g3810bn-tl|\
netgear,dg834g-v4|\
usrobotics,usr9108)
	migrate_leds "^.*::="
	;;
esac

remove_devicename_leds

migrations_apply system

exit 0
