# package/lean/luci-app-travelmate-v2/root/etc/uci-defaults/40_luci-travelmate
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@travelmate[-1]
	add ucitrack travelmate
	set ucitrack.@travelmate[-1].init=travelmate
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
