# package/others/luci-theme-darkmatter/root/etc/uci-defaults/30_luci-theme-darkmatter
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Darkmatter=/luci-static/darkmatter
	commit luci
EOF
exit 0
