# package/Twinzo/luci-app-serverchan/root/etc/init.d/serverchan
#!/bin/sh /etc/rc.common

START=99
STOP=10

start() {	
	/usr/bin/serverchan/serverchan &
	echo "serverchan is starting now ..."
}

stop() {
	kill -9 `pgrep -f "/usr/bin/serverchan/serverchan"` 2>/dev/null
	echo "serverchan exit ..."
}

restart(){
	stop
	sleep 1
	start
} 
