# package/lean/RNDIS/root/etc/uci-defaults/RNDIS
#!/bin/sh 
#
#----------------------------------------------
#Copyright (C) 2019 Twinzo1 <1282055288@qq.com>
#----------------------------------------------
#
uci -q batch <<-EOF >/dev/null
	set network.usbwan="interface"
	set network.usbwan.ifname="usb0"
	set network.usbwan.proto="dhcp"
	commit network
    EOF
exit 0
