# package/Twinzo/luci-app-idns/root/etc/uci-defaults/luci-idns
#!/bin/sh

chmod 777 /etc/init.d/idns
rm -f /tmp/luci-indexcache
exit 0
