# package/Twinzo/luci-app-idns/root/etc/uci-defaults/luci-idns
#!/bin/sh

rm -f /tmp/luci-indexcache
chmod 777 /etc/init.d/idns
exit 0