# target/linux/ath79/base-files/etc/hotplug.d/ieee80211/00-wmac-migration
#!/bin/sh

WMAC_PATH_CHANGED=0

. /lib/functions.sh

migrate_wmac_path() {
	local section="$1"
	local path

	config_get path ${section} path
	case ${path} in
		"platform/qca955x_wmac")
			path="platform/ahb/ahb:apb/18100000.wmac"
			WMAC_PATH_CHANGED=1
		;;
		"platform/ar933x_wmac")
			path="platform/ahb/18100000.wmac"
			WMAC_PATH_CHANGED=1
		;;
		*)
			return 0
		;;
	esac

	uci set wireless.${section}.path=${path}
}

[ "${ACTION}" = "add" ] && {
	[ ! -e /etc/config/wireless ] && return 0

	config_load wireless
	config_foreach migrate_wmac_path wifi-device

	[ "${WMAC_PATH_CHANGED}" = "1" ] && uci commit wireless
}
