# package/Twinzo/luci-app-watchpig/root/etc/uci-defaults/luci-watchpig
#!/bin/sh 

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@watchpig[-1]
	add ucitrack watchpig
    	set ucitrack.@watchpig[-1].init=watchpig
    	commit ucitrack
EOF

uci -q set watchpig.@mac[0].interface="wan"
[ -n $`uci -q get network.wan_dev` ] && uci -q set watchpig.@mac[0].interface="wan_dev"
uci -q set watchpig.@mac[0].macaddr="$(uci -q get network.wan.macaddr||uci -q get network.wan_dev.macaddr)"
uci commit watchpig

chmod 755 /etc/init.d/watchpig
chmod 755 /usr/bin/watchpig.sh

rm -f /tmp/luci-indexcache
exit 0
