# package/lean/luci-app-wrtbwmon/root/etc/uci-defaults/luci-wrtbwmon
#!/bin/sh

/etc/init.d/wrtbwmon enable
/etc/init.d/wrtbwmon start
exit 0
