# package/community/openwrt-package/luci-app-nginx-pingos/root/etc/uci-defaults/luci-app-nginx-pingos
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@pingos[-1]
	add ucitrack pingos
	set ucitrack.@pingos[-1].init=pingos
	commit ucitrack
EOF

rm -rf /tmp/luci-*cache
exit 0
