# customfeeds/packages/utils/bash/files/etc/profile.d/sys_bashrc.sh
[ -z "$BASH" ] || [ "$BASH" = /bin/sh ] || [ \! -s /etc/bash.bashrc ] || . /etc/bash.bashrc
