# customfeeds/packages/net/apinger/files/apinger.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=80

SERVICE_USE_PID=1

start() {
	service_start /usr/sbin/apinger
}

stop() {
	service_stop /usr/sbin/apinger
}

reload() {
	service_reload /usr/sbin/apinger
}
