# package/lean/luci-app-cifsd/root/etc/uci-defaults/cifsd
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@ksmbd[-1]
	add ucitrack ksmbd
	set ucitrack.@ksmbd[-1].init=ksmbd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
