# package/lean/luci-app-cifsd/root/etc/uci-defaults/cifsd
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@cifsd[-1]
	add ucitrack cifsd
	set ucitrack.@cifsd[-1].init=cifsd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
