# target/linux/mediatek/mt7623/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
::askconsole:/usr/libexec/login.sh
tty1::askfirst:/usr/libexec/login.sh
tty2::askfirst:/usr/libexec/login.sh
tty3::askfirst:/usr/libexec/login.sh
tty4::askfirst:/usr/libexec/login.sh
tty5::askfirst:/usr/libexec/login.sh
tty6::askfirst:/usr/libexec/login.sh
