# customfeeds/luci/applications/luci-app-ddns/root/etc/uci-defaults/40_luci-ddns
#!/bin/sh

# no longer needed for "Save and Apply" to restart ddns
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@ddns[-1]
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
return 0
