# customfeeds/packages/net/fakeidentd/files/fakeidentd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2012 OpenWrt.org

START=85

start() {
	service_start /usr/sbin/fakeidentd -r
}

stop() {
	service_stop /usr/sbin/fakeidentd
}
