# package/my/luci-app-flowoffload/root/etc/hotplug.d/iface/55-fullconenat
#!/bin/sh
# Copyright (C) 2007 OpenWrt.org

nat1_enable=$(uci get flowoffload.@flow[0].nat1_enabled 2>/dev/null)

/etc/init.d/flowoffload enabled && {

	[ "$ACTION" = ifup -o "$ACTION" = wan ] && {
		if [ $nat1_enable -eq 1 ]; then
		/etc/init.d/flowoffload restart
		fi
	}
	
	[ "$ACTION" = ifdown -o "$ACTION" = wan ] && {
		if [ $nat1_enable -eq 1 ]; then
		/etc/init.d/flowoffload restart
		fi
	}

	[ "$ACTION" = ifup -o "$ACTION" = lan ] && {
		if [ $nat1_enable -eq 1 ]; then
		/etc/init.d/flowoffload restart
		fi
	}
	
	[ "$ACTION" = ifdown -o "$ACTION" = lan ] && {
		if [ $nat1_enable -eq 1 ]; then
		/etc/init.d/flowoffload restart
		fi
	}
	
}
