# customfeeds/luci/applications/luci-app-vnstat/root/etc/uci-defaults/40_luci-vnstat
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@vnstat[-1]
	add ucitrack vnstat
	set ucitrack.@vnstat[-1].init=vnstat
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
