# customfeeds/packages/net/rsync/files/rsyncd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2011-2014 OpenWrt.org

START=90
STOP=10

USE_PROCD=1
PROG=/usr/bin/rsync

start_service() {
	procd_open_instance
	procd_set_param command "$PROG" --daemon --no-detach
	procd_close_instance
}
