# customfeeds/packages/net/krb5/files/krb5kdc
#!/bin/sh /etc/rc.common
# Copyright (C) 2009 OpenWrt.org

START=60

start() {
	mkdir -p /var/tmp
	mkdir -p /etc/krb5kdc
	
	[ -f /etc/krb5kdc/principal ] || ( echo; echo ) | kdb5_util create -s
	
	/usr/sbin/krb5kdc
	[ -x /usr/sbin/kadmind ] && /usr/sbin/kadmind
}

stop() {
	killall krb5kdc 2> /dev/null
	[ -x /usr/sbin/kadmind ] && killall kadmind 2> /dev/null
}
