# package/lean/luci-app-jd-dailybonus/root/etc/uci-defaults/luci-jd-dailybonus
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@jd-dailybonus[-1]
	add ucitrack jd-dailybonus
	set ucitrack.@jd-dailybonus[-1].init=jd-dailybonus
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
