# package/lean/luci-app-arpbind/root/etc/hotplug.d/iface/50-arpbind
#!/bin/sh
[ "$ACTION" = "ifup" ] && /etc/init.d/arpbind start
