# customfeeds/luci/applications/luci-app-shairplay/root/etc/uci-defaults/40_luci-shairplay
#!/bin/sh

# needed for "Save and Apply" to restart shairplay
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@shairplay[-1]
	add ucitrack shairplay
	set ucitrack.@shairplay[-1].init="shairplay"
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
