# package/my/luci-app-smartdns/root/etc/uci-defaults/50_luci-smartdns
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@smartdns[-1]
	add ucitrack smartdns
	set ucitrack.@smartdns[-1].init=smartdns
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
