# package/lean/luci-app-unblockmusic/root/etc/hotplug.d/iface/099-unblockmusic
#!/bin/sh

[ "$ACTION" = ifup ] || exit 0

sleep 10
/etc/init.d/unblockmusic restart
