# package/Twinzo/luci-app-oafControl/root/etc/init.d/oafControl
#!/bin/sh /etc/rc.common

#----------------------------------------------
#Copyright (C) 2019 Twinzo1 <1282055288@qq.com>
#----------------------------------------------
START=98

start(){
	global_enable=$(uci get oafControl.@global[0].global_enabled 2>/dev/null)
	global_run=$(uci get oafControl.@global[0].global_run 2>/dev/null)
	auto_update=$(uci get oafControl.@global[0].global_update 2>/dev/null)
	[ ! -d "/etc/oafControl" ] && mkdir /etc/oafControl
	[ ! -d "/tmp/oafControl" ] && mkdir /tmp/oafControl
	
	if [ "$global_run" != "0" ]; then
		sed -i '/#remote_control/d' /etc/crontabs/root
		echo "*/${global_run} * * * * /usr/bin/oafControl.sh reomte_control #remote_control" >> /etc/crontabs/root
	else
		sed -i '/#remote_control/d' /etc/crontabs/root
	fi
	
	if [ "$auto_update" != "0" ]; then
		sed -i '/#oafControl_timing_update/d' /etc/crontabs/root
		echo "30 3 * * * /usr/bin/oafControl.sh auto_update #oafControl_timing_update" >> /etc/crontabs/root
		echo "`date '+%A %Y-%m-%d %H:%M:%S'`：每天凌晨三点半自动更新特征库" >> /tmp/oafControl/oafControl.log
	else
		sed -i '/#oafControl_timing_update/d' /etc/crontabs/root
	fi
	
	sed -i '/#global_control/d' /etc/crontabs/root
	echo "00 3 * * * /usr/bin/oafControl.sh global_control #global_control" >> /etc/crontabs/root
	echo "`date '+%A %Y-%m-%d %H:%M:%S'`：每天凌晨三点检测当天加班还是放假" >> /tmp/oafControl/oafControl.log
	
	if [ "$global_enable"x != "1"x ]; then
		/etc/init.d/oafControl disable
		/usr/bin/oafControl.sh global_stop all
	else 
		/etc/init.d/oafControl enable
		/usr/bin/oafControl.sh global_start
		/usr/bin/oafControl.sh init_weekday
		/usr/bin/oafControl.sh init_workday
		/usr/bin/oafControl.sh init_vacation
		/usr/bin/oafControl.sh init_overtime
	fi
}

stop(){
	echo "`date '+%A %Y-%m-%d %H:%M:%S'`：oafControl重启" >> /tmp/oafControl/oafControl.log
}