# package/diy/luci-app-cpe/root/etc/init.d/cpeinit
#!/bin/sh /etc/rc.common

START=99

start() {
  /bin/sh /usr/share/cpe/rssi &
}