# package/diy/luci-app-fan/root/etc/uci-defaults/luci-fan
#!/bin/sh

/usr/libexec/fan-control init

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@luci-fan[-1]
	add ucitrack luci-fan
	set ucitrack.@luci-fan[-1].init=luci-fan
	commit ucitrack
EOF

exit 0
