# target/linux/ipq60xx/base-files/etc/hotplug.d/firmware/11-ath11k-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

case "$FIRMWARE" in
"ath11k/IPQ6018/hw1.0/cal-ahb-c000000.wifi.bin")
	case $(board_name) in
	linksys,mr7350)
		caldata_extract "0:art" 0x1000 0x10000
		;;
	tplink,eap610-outdoor)
		# /tmp/factory_data should be mounted by preinit
		cp /tmp/factory_data/radio "/lib/firmware/$FIRMWARE"
		;;
	*)
		caldata_die "Don't know how to read caldata for $(board_name)"
                ;;
	esac
	;;
esac
