# target/linux/ipq60xx/base-files/etc/hotplug.d/firmware/11-ath11k-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

mount_factory_data() {
	local firmware_name=$1
	local part
	. /lib/functions/system.sh


	part=$(find_mtd_chardev "factory_data")
	ubiattach -p "$part"
	mkdir /tmp/facdata
	mount -t ubifs,ro ubi0:ubi_factory_data /tmp/facdata
	cp /tmp/facdata/radio "/lib/firmware/$firmware_name"
}

case "$FIRMWARE" in
"ath11k/IPQ6018/hw1.0/cal-ahb-c000000.wifi.bin")
	mount_factory_data $FIRMWARE
	;;
*)
	exit 1
	;;
esac
