# package/base-files/files/etc/init.d/fa-wifiap
#!/bin/sh /etc/rc.common

START=11
STOP=15
 
start() {
	echo "fa-wifiap started"
	/bin/sh /usr/bin/fix_wifi_ap.sh > /dev/null&
}

stop() {
	echo "fa-wifiap: do nothing"
}
