# package/base-files/files/etc/init.d/sysctl
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=11

set_vm_min_free() {
	mem="$(grep MemTotal /proc/meminfo  | awk '{print $2}')"
	if [ "$mem" -gt 65536 ]; then # 128M
		val=16384
	elif [ "$mem" -gt 32768 ]; then # 64M
		val=8192
	elif [ "$mem" -gt 16384 ]; then # 32M
		val=1024
	else
		return
	fi
	sysctl -qw vm.min_free_kbytes="$val"
}

start() {
	set_vm_min_free
	for CONF in /etc/sysctl.conf /etc/sysctl.d/*.conf; do
		[ -f "$CONF" ] && sysctl -p "$CONF" -e >&-
	done
}
