# package/base-files/files/etc/init.d/umount
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

STOP=99
stop() {
	sync
	/bin/umount -a -d -r
}
