# target/linux/kirkwood/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	linksys,audi|\
	linksys,viper)
		mtd resetbc s_env || true
		;;
	esac
}
