# target/linux/kirkwood/base-files/etc/init.d/hwmon_fancontrol
#!/bin/sh /etc/rc.common

START=98

boot() {
	# configuring (lm85/lm63) onboard temp/fan controller to run the fan on its own
	# for more information, please read https://www.kernel.org/doc/Documentation/hwmon/sysfs-interface

	case $(board_name) in
	iom,ix2-200)
		path_to_hwmon='/sys/class/hwmon/hwmon0'
		echo 2 > "$path_to_hwmon/pwm1_enable" # fan is on pwm1
		;;
	seagate,blackarmor-nas220)
		path_to_hwmon='/sys/devices/platform/ocp@f1000000/f1011000.i2c/i2c-0/0-002e/hwmon/hwmon0'
		# adt7476 fan control chip. 3 temp sensors. Set to 1/4 speed at 35C and max speed at 48C.
		echo 7 > "$path_to_hwmon/pwm1_auto_channels_temp"
		echo 64 > "$path_to_hwmon/pwm1_auto_point1_pwm"
		echo 255 > "$path_to_hwmon/pwm1_auto_point2_pwm"
		echo 35000 > "$path_to_hwmon/temp1_auto_point1_temp"
		echo 48000 > "$path_to_hwmon/temp1_auto_point2_temp"
		echo 35000 > "$path_to_hwmon/temp2_auto_point1_temp"
		echo 48000 > "$path_to_hwmon/temp2_auto_point2_temp"
		echo 35000 > "$path_to_hwmon/temp3_auto_point1_temp"
		echo 48000 > "$path_to_hwmon/temp3_auto_point2_temp"
		echo 2 > "$path_to_hwmon/pwm1_enable"
		;;
	zyxel,nsa310b)
		path_to_hwmon='/sys/devices/platform/ocp@f1000000/f1011000.i2c/i2c-0/0-002e/hwmon/hwmon0'
		echo 2 > "$path_to_hwmon/pwm1_enable" # fan is on pwm1
		echo 1 > "$path_to_hwmon/pwm1_auto_channels" # temp1 is the only one that changes
		echo 23000 > "$path_to_hwmon/temp1_auto_temp_min"
		echo 43000 > "$path_to_hwmon/temp1_auto_temp_max" # next step is 49600 millicelsius, or 50 celsius, 43 celsius is better
		;;
	esac
}
