# target/linux/kirkwood/base-files/etc/init.d/hwmon_fancontrol
#!/bin/sh /etc/rc.common

START=98
boot() {
. /lib/functions.sh

#configuring (lm85/lm63) onboard temp/fan controller to run the fan on its own
#for more information, please read https://www.kernel.org/doc/Documentation/hwmon/sysfs-interface

case $(board_name) in
zyxel,nsa310b)
	path_to_hwmon='/sys/devices/platform/ocp@f1000000/f1011000.i2c/i2c-0/0-002e/hwmon/hwmon0'
	echo 2 > "$path_to_hwmon/pwm1_enable" # fan is on pwm1
	echo 1 > "$path_to_hwmon/pwm1_auto_channels" # temp1 is the only one that changes
	echo 23000 > "$path_to_hwmon/temp1_auto_temp_min"
	echo 43000 > "$path_to_hwmon/temp1_auto_temp_max" # next step is 49600 millicelsius, or 50 celsius, 43 celsius is better
	;;
iom,ix2-200)
	path_to_hwmon='/sys/class/hwmon/hwmon0'
	echo 2 > "$path_to_hwmon/pwm1_enable" # fan is on pwm1
	;;
esac
}
