# target/linux/kirkwood/base-files/etc/init.d/linksys_recovery
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 OpenWrt.org

START=97
boot() {
. /lib/functions.sh

case $(board_name) in
	linksys,audi|linksys,viper)
		# make sure auto_recovery in uboot is always on
		AUTO_RECOVERY_ENA="`fw_printenv -n auto_recovery`"
		if [ "$AUTO_RECOVERY_ENA" != "yes" ] ; then
			fw_setenv auto_recovery yes
		fi
		# reset the boot counter
		mtd resetbc s_env
		;;
esac
}
