# target/linux/rockchip/rk35xx/base-files/etc/uci-defaults/13_opkg_update
#!/bin/sh

sed -i 's/rk35xx/armv8/g' /etc/opkg/distfeeds.conf

exit 0

