# target/linux/mvebu/cortexa9/base-files/etc/uci-defaults/05_fix-compat-version
. /lib/functions.sh

case "$(board_name)" in
	linksys,wrt1900ac-v1|\
        linksys,wrt32x)
		uci set system.@system[0].compat_version="2.0"
		uci commit system
		;;
esac

exit 0
