# customfeeds/luci/applications/luci-app-minidlna/root/etc/uci-defaults/40_luci-minidlna
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.minidlna
	set ucitrack.minidlna=minidlna
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
