# customfeeds/packages/net/dhcpcd/files/dhcpcd.init
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2014 OpenWrt.org
#
# This is free software, licensed under the GNU General Public License v2.
# See /LICENSE for more information.
#

START=30
STOP=85
USE_PROCD=1

start_service()
{
	procd_open_instance
	procd_set_param command /sbin/dhcpcd -B
	procd_set_param respawn
	procd_close_instance
}

reload_service()
{
	/sbin/dhcpcd -n
}
