# customfeeds/luci/applications/luci-app-p910nd/root/etc/uci-defaults/40_luci-p910nd
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@p910nd[-1]
	add ucitrack p910nd
	set ucitrack.@p910nd[-1].init=p910nd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
