# package/Twinzo/luci-app-chinadns/root/etc/uci-defaults/luci-chinadns
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@chinadns[-1]
	add ucitrack chinadns
	set ucitrack.@chinadns[-1].init=chinadns
	commit ucitrack
EOF

chmod 777 /etc/init.d/chinadns
chmod 777 /usr/lib/chinadns/update
rm -f /tmp/luci-indexcache
exit 0
