# customfeeds/luci/applications/luci-app-bcp38/root/etc/uci-defaults/60_luci-bcp38
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@bcp38[-1]
	add ucitrack bcp38
        add_list ucitrack.@bcp38[0].affects=firewall
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
