# package/lean/luci-app-apfree_wifidog/root/etc/uci-defaults/90_luci_wifidog
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@wifidog[-1]
	add ucitrack wifidog
	set ucitrack.@wifidog[-1].init=wifidog
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0