# target/linux/mediatek/base-files/etc/uci-defaults/99-net-ps
uci set network.globals.default_rps_val=14
uci set network.globals.default_rps_flow_cnt=256
uci set network.globals.default_xps_val=14
uci set network.globals.default_ps=1
uci set network.eth0=device
uci set network.eth0.name=eth0
uci set network.lan0=device
uci set network.lan0.name=lan0
uci set network.lan1=device
uci set network.lan1.name=lan1
uci set network.lan2=device
uci set network.lan2.name=lan2
uci set network.lan3=device
uci set network.lan3.name=lan3
uci commit
exit 0
