# target/linux/ipq807x/base-files/etc/hotplug.d/firmware/11-ath10k-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"ath10k/cal-pci-0000:01:00.0.bin")
	case "$board" in
	xiaomi,ax3600)
		caldata_extract "0:art" 0x33000 0x844
		;;
	esac
	;;
*)
	exit 1
	;;
esac
