# target/linux/ipq807x/base-files/etc/hotplug.d/firmware/11-ath11k-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"ath11k/IPQ8074/hw2.0/cal-ahb-c000000.wifi.bin")
	case "$board" in
	qnap,301w)
		caldata_extract "0:art" 0x1000 0x20000
		;;
	zte,mf269)
		caldata_extract "0:ART" 0x1000 0x20000
		;;
	esac
	;;
*)
	exit 1
	;;
esac
