# target/linux/ipq807x/base-files/etc/hotplug.d/firmware/11-ath11k-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"ath11k/IPQ8074/hw2.0/caldata.bin")
	case "$board" in
	netgear,sxr80|\
	redmi,ax6|\
	xiaomi,ax3600|\
	xiaomi,ax9000)
		caldata_extract "0:art" 0x1000 0x20000
		;;
	zte,mf269)
		caldata_extract "0:ART" 0x1000 0x20000
		;;
	esac
	;;
"ath11k/QCN9074/hw1.0/caldata.bin")
	case "$board" in
	xiaomi,ax9000)
		caldata_extract "0:art" 0x26800 0x20000
		;;
	esac
	;;
*)
	exit 1
	;;
esac
