# target/linux/ipq807x/base-files/etc/hotplug.d/firmware/11-ath11k-caldata
#!/bin/sh -x

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"IPQ8074/caldata.bin")
	case "$board" in
	xiaomi,ax3600)
		caldata_extract "0:ART" 0x1000 0x20000
		;;
	esac
	;;
*)
	exit 1
	;;
esac
