# customfeeds/packages/net/bcp38/files/bcp38.init
#!/bin/sh /etc/rc.common

START=20

USE_PROCD=1
NAME=bcp38

service_triggers()
{
	procd_add_config_trigger "config.change" "bcp38" /etc/init.d/firewall reload
}
