# package/lean/luci-theme-atmaterial/root/etc/uci-defaults/30-luci-theme-atmaterial
#!/bin/sh
uci batch <<-EOF
	set luci.themes.AdvancedTomatoMaterial=/luci-static/atmaterial
	set luci.main.mediaurlbase=/luci-static/atmaterial
	commit luci
EOF
exit 0
