# package/lean/luci-app-aliddns/root/etc/uci-defaults/luci-app-aliddns
#!/bin/sh


rm -f /tmp/luci-indexcache
exit 0
