# package/lean/luci-app-aliddns/root/etc/uci-defaults/luci-aliddns
#!/bin/sh 

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@aliddns[-1]
	add ucitrack aliddns
	set ucitrack.@aliddns[-1].init=aliddns
	commit ucitrack
EOF

exit 0
