# target/linux/brcm63xx/base-files/etc/uci-defaults/09_fix_crc
#!/bin/sh
#
# Copyright (C) 2007 OpenWrt.org
#
#

. /lib/functions.sh

do_fixcrc() {
	mtd fixtrx linux
}

case "$(board_name)" in
	a4001n |\
	a4001n1 |\
	ar-5315u |\
	ar-5381u |\
	ar-5387un |\
	bcm96328avng |\
	bcm963281tan |\
	cpva502p |\
	cpva642 |\
	ct-6373 |\
	dsl-274xb-f |\
	evg2000 |\
	hg622 |\
	magic |\
	p870hw-51a_v2 |\
	r1000h |\
	r5010un_v2 |\
	rta770bw |\
	rta770w |\
	spw303v |\
	v2110 |\
	v2500v_bb |\
	vh4032n |\
	vr-3025u |\
	vr-3025un |\
	vr-3026e |\
	wap-5813n)
		do_fixcrc
		;;
esac
