# customfeeds/packages/net/switchdev-poller/files/switchdev-poller.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2018-2019 Luiz Angelo Daros de Luca <luizluca@gmail.com>

START=99
USE_PROCD=1

start_service() {
	procd_open_instance switchdev-poller
	procd_set_param command /usr/lib/switchdev-poller/switchdev-poller
	procd_set_param stdout 1
	procd_set_param stderr 0 # enable for extra debug info
	procd_set_param respawn 3600 5 5
	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger network
}
