# package/lean/luci-app-amule/root/etc/hotplug.d/iface/20-amule
#!/bin/sh
# Copyright (C) 2007 OpenWrt.org

/etc/init.d/amule enabled && {

	[ "$ACTION" = "ifup" ] && {
		/etc/init.d/amule start
	}

	[ "$ACTION" = "ifdown" ] && {
		/etc/init.d/amule stop
	}
	
}
