# package/lean/luci-app-baidupcs-web/root/etc/uci-defaults/luci-baidupcs-web
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@baidupcs-web[-1]
	add ucitrack baidupcs-web
	set ucitrack.@baidupcs-web[-1].init=baidupcs-web
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
