# customfeeds/packages/net/kplex/files/kplex.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 OpenWrt.org

START=99

USE_PROCD=1
PROG=/usr/bin/kplex

CONF_FILE=/etc/kplex.conf

start_service() {
	[ -e "$CONF_FILE" ] || return 1
	
	procd_open_instance
	procd_set_param command "$PROG" -f "$CONF_FILE"
	procd_set_param respawn	
	procd_set_param file "$CONF_FILE"
	procd_close_instance
}
