# customfeeds/packages/net/horst/horst.init
#!/bin/sh /etc/rc.common

START=94

USE_PROCD=1
PROG=/usr/sbin/horst

start_service() {
	config_load horst
	config_get dis "init" 'disabled' '0'
	[ "$dis" == "0" ] || exit

	config_get intf "init" 'interface' 'wlan0'
	config_get scan "init" 'scan' '0'
	[ "$scan" == "0" ] && scan='' || scan="-s"

	procd_open_instance
	procd_set_param command $PROG -q -i $intf $scan -N
	procd_set_param netdev $intf
	procd_set_param respawn
	procd_close_instance
}
