# package/lean/luci-app-filetransfer/etc/uci-defaults/luci-filetransfer
#!/bin/sh

sed -i 's/cbi.submit\"] = true/cbi.submit\"] = \"1\"/g' /usr/lib/lua/luci/dispatcher.lua

rm -f /tmp/luci-indexcache
exit 0
