# customfeeds/packages/net/net-snmp/files/snmptrapd.init
#!/bin/sh /etc/rc.common

START=50

USE_PROCD=1
PROG="/usr/sbin/snmptrapd"

start_service() {
	procd_open_instance

	procd_set_param command $PROG -Lf /dev/null -f
	procd_set_param respawn

	procd_close_instance
}
