# customfeeds/luci/applications/luci-app-radicale/root/etc/uci-defaults/40_luci-radicale
#!/bin/sh

# no longer needed for "Save and Apply" to restart radicale
# luci-app-radicale calls /etc/init.d/radicale reload
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@radicale[-1]
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache

exit 0
