# customfeeds/packages/net/softethervpn/files/vpnbridge.init
#!/bin/sh /etc/rc.common
START=60

start(){
        logger -t 'softethervpn' "Starting softether vpnbridge service."
	/usr/bin/env LANG=en_US.UTF-8 /usr/libexec/softethervpn/vpnbridge start
}

stop(){
        logger -t 'softethervpn' "Stopping softether vpnbridge service."
	/usr/bin/env LANG=en_US.UTF-8 /usr/libexec/softethervpn/vpnbridge stop
}
