# customfeeds/packages/net/softethervpn/files/vpnclient.init
#!/bin/sh /etc/rc.common
START=60

start(){
        logger -t 'softethervpn' "Starting softether vpnclient service."
	/usr/bin/env LANG=en_US.UTF-8 /usr/libexec/softethervpn/vpnclient start
}

stop(){
        logger -t 'softethervpn' "Stopping softether vpnclient service."
	/usr/bin/env LANG=en_US.UTF-8 /usr/libexec/softethervpn/vpnclient stop
}
