# customfeeds/packages/net/opentracker/files/opentracker.init
#!/bin/sh /etc/rc.common

START=10
STOP=15

NAME="opentracker"
PROG="/usr/bin/opentracker"
OPTIONS="-f /etc/opentracker.conf"
USE_PROCD=1

start_service()
{
	procd_open_instance
	procd_set_param command $PROG $OPTIONS
	procd_close_instance
}
