# customfeeds/packages/net/opentracker/files/opentracker6.init
#!/bin/sh /etc/rc.common

START=10
STOP=15

NAME="opentracker6"
PROG="/usr/bin/opentracker6"
OPTIONS="-f /etc/opentracker6.conf"
USE_PROCD=1

start_service()
{
	procd_open_instance
	procd_set_param command $PROG $OPTIONS
	procd_close_instance
}
