# package/community/openwrt-package/luci-app-socat/root/etc/uci-defaults/luci-app-socat
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete firewall.socat
	set firewall.socat=include
	set firewall.socat.type=script
	set firewall.socat.path=/var/etc/socat.include
	set firewall.socat.reload=1
EOF

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@socat[-1]
	add ucitrack socat
	set ucitrack.@socat[-1].init=socat
	commit ucitrack
EOF

rm -rf /tmp/luci-*cache
exit 0
