# customfeeds/packages/utils/cgroupfs-mount/files/cgroupfs-mount.init
#!/bin/sh /etc/rc.common

START=01

boot() {
	# Procd mounts non-hierarchical cgroupfs so unmount first before cgroupfs-mount
	if mountpoint -q /sys/fs/cgroup; then
		umount /sys/fs/cgroup/
	fi

	cgroupfs-mount
}
