# package/lean/luci-app-xlnetacc/root/etc/uci-defaults/luci-xlnetacc
uci batch <<EOF
	delete ucitrack.@xlnetacc[-1]
	add ucitrack xlnetacc
	set ucitrack.@xlnetacc[-1].init=xlnetacc
	commit ucitrack
EOF
