# package/lean/luci-app-familycloud/root/etc/uci-defaults/familycloud
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@familycloud[-1]
	add ucitrack familycloud
	set ucitrack.@familycloud[-1].init=familycloud
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
