# package/lean/luci-app-fogvdn/root/etc/uci-defaults/30-run_post_update_cmd
#!/bin/sh
mkdir /run
/etc/pear/pear_update/post_command.sh
