# package/lean/luci-app-fogvdn/root/etc/uci-defaults/luci-fogvdn
#!/bin/sh
rm -f /tmp/luci-indexcache
#uci track
cfg_id=$(uci add ucitrack fogvdn)
uci set ucitrack.$cfg_id.init="fogvdn"
uci commit ucitrack
exit 0
