# customfeeds/packages/libs/avahi/files/avahi-daemon.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=61

USE_PROCD=1
PROG=avahi-daemon

start_service() {
	procd_open_instance
	procd_set_param command "$PROG"
	procd_append_param command -s
	procd_set_param respawn
	procd_close_instance
}

reload_service() {
	procd_send_signal "$PROG"
}
