# package/system/fstools/files/fstab.init
#!/bin/sh /etc/rc.common
# (C) 2013 openwrt.org

START=40

boot() {
	/sbin/block mount
}

start() {
	echo "this file has been obsoleted. please call \"/sbin/block mount\" directly"
}

stop() {
	/sbin/block umount
}
