# customfeeds/luci/applications/luci-app-mjpg-streamer/root/etc/uci-defaults/40_luci-mjpg-streamer
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@mjpg-streamer[-1]
	add ucitrack mjpg-streamer
	set ucitrack.@mjpg-streamer[-1].init=mjpg-streamer
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
