# customfeeds/packages/net/dnsdist/files/dnsdist.init
#!/bin/sh /etc/rc.common
START=99

USE_PROCD=1

start_service() {
	config_load dnsdist
	local cfg=general
	local enabled

	config_get_bool enabled "$cfg" 'enabled' 1

	[ $enabled -gt 0 ] || return 1

	procd_open_instance
	procd_set_param command dnsdist --supervised -C /etc/dnsdist.conf
	procd_set_param file /etc/dnsdist.conf
	procd_set_param respawn
	procd_close_instance
}
