# target/linux/rockchip-rk3328/base-files/etc/hotplug.d/iface/12-disable-rk3328-eth-offloading
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh

[ -n "$INTERFACE" ] || exit 3
[ "$INTERFACE" == "wan" ] || exit 4
[ "$ACTION" == "ifup" ] || exit 1

BOARD=$(cat /tmp/sysinfo/board_name | cut -d , -f2)
logger -t disable-offloading "board = $BOARD"
[ "$BOARD" == "nanopi-r2" -o "$BOARD" == "nanopi-r2s" ] || exit 2
logger -t disable-offloading "rk3328 board detected"

ETHTOOL=/usr/sbin/ethtool
[ -f $ETHTOOL ] || exit 5

[ -d /sys/devices/platform/ff540000.ethernet/net/eth0 ] || exit 6
$ETHTOOL -K eth0 rx off tx off

logger -t disable-offloading "disabed rk3328 ethernet tcp/udp offloading tx/rx"
exit 0
