# customfeeds/packages/utils/triggerhappy/files/triggerhappy.init
#!/bin/sh /etc/rc.common
START=93

start() {
	/usr/sbin/thd --socket /tmp/triggerhappy.socket --triggers /etc/triggerhappy/triggers.d/ --daemon /dev/input/event*
}

stop() {
	/usr/sbin/th-cmd --socket /tmp/triggerhappy.socket --quit
}
