# target/linux/rockchip/armv8/base-files/etc/hotplug.d/usb/15-usb-wifi-config
#!/bin/sh

if [ "${ACTION}" = "add" ]; then
  sleep 5
  [ -n "$(cat /proc/net/wireless | grep wlan)" ] && wifi config && wifi
fi
