# target/linux/rockchip/armv8/base-files/etc/hotplug.d/usb/50-aic8800_usb
[ "${DRIVER}" = "aic8800_usb" ] && {
  sleep 3
  if [ ! -f /etc/config/qmodem ]; then
    /sbin/wifi config && /sbin/wifi
  else
    /sbin/wifi
  fi
}
