# package/kernel/lantiq/ltq-vmmc/files/vmmc.init
#!/bin/sh /etc/rc.common
#
# Activate Voice CPE TAPI subsystem LL driver for VMMC

START=31

start() {
	[ ! -c /dev/vmmc10 ] && {
		mknod /dev/vmmc10 c 122 10
		mknod /dev/vmmc11 c 122 11
		mknod /dev/vmmc12 c 122 12
		mknod /dev/vmmc13 c 122 13
		mknod /dev/vmmc14 c 122 14
		mknod /dev/vmmc15 c 122 15
		mknod /dev/vmmc16 c 122 16
		mknod /dev/vmmc17 c 122 17
		mknod /dev/vmmc18 c 122 18
	}
}
