# package/lean/luci-theme-argon/root/etc/uci-defaults/30_luci-theme-argon
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Argon=/luci-static/argon
	set luci.main.mediaurlbase=/luci-static/argon
	commit luci
EOF
exit 0
