# package/lean/luci-theme-argon/root/etc/uci-defaults/30_luci-theme-argon
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Argon=/luci-static/argon
	set luci.main.mediaurlbase=/luci-static/argon
	commit luci
EOF
[ -f /usr/lib/lua/luci/view/themes/argon/out_header_login.htm ] && mv -f /usr/lib/lua/luci/view/themes/argon/out_header_login.htm /usr/lib/lua/luci/view/header_login.htm
rm -Rf /var/luci-modulecache
rm -Rf /var/luci-indexcache
exit 0
