# customfeeds/packages/net/mwan3/files/etc/uci-defaults/mwan3-migrate-flush_conntrack
#!/bin/sh

. /lib/functions.sh

mwan3_migrate_flush_conntrack() {
	local iface="$1"

	config_get value "${iface}" flush_conntrack
	case $value in
		always)
			uci_remove mwan3 "$iface" flush_conntrack
			uci_add_list mwan3 "$iface" flush_conntrack ifup
			uci_add_list mwan3 "$iface" flush_conntrack ifdown
			;;
		never)
			uci_remove mwan3 "$iface" flush_conntrack
			;;
	esac

	uci_commit mwan3
}

config_load mwan3
config_foreach mwan3_migrate_flush_conntrack interface

exit 0
