# target/linux/rockchip/armv8/base-files/etc/init.d/fa-fancontrol
#!/bin/sh /etc/rc.common

START=21
start() {
    /usr/bin/fa-fancontrol.sh &
}

