# target/linux/rockchip/armv8/base-files/etc/init.d/startup
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=99

. /lib/functions/system.sh

boot() {
	case $(board_name) in
	rongpin,king3399)
		echo 50 > /sys/class/gpio/export && echo high > /sys/class/gpio/gpio50/direction
        echo 56 > /sys/class/gpio/export && echo high > /sys/class/gpio/gpio56/direction
		;;
	esac
}