# package/lienol/luci-app-control-weburl/root/etc/uci-defaults/luci-app-control-weburl
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@weburl[-1]
	add ucitrack weburl
	set ucitrack.@weburl[-1].init=weburl
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
