# package/lean/luci-app-kuainiao/root/etc/uci-defaults/luci-kuainiao
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@kuainiao[-1]
	add ucitrack kuainiao
	set ucitrack.@kuainiao[-1].init=kuainiao
	commit ucitrack
EOF
exit 0