# target/linux/orion/generic/base-files/etc/uci-defaults/09_hardware
#!/bin/sh
#
# Copyright (C) 2010 OpenWrt.org
#
# This is free software, licensed under the GNU General Public License v2.
# See /LICENSE for more information.
#

#
# This script sets system defaults for the hardware on firstboot
#

. /lib/functions.sh

wrt350nv2_default() {
# leds
	uci batch <<__EOF
set system.power_led=led
set system.power_led.name='Power LED (green)'
set system.power_led.sysfs='wrt350nv2:green:power'
set system.power_led.default='1'
set system.wifi_led=led
set system.wifi_led.name='Wireless LED (green)'
set system.wifi_led.sysfs='wrt350nv2:green:wireless'
set system.wifi_led.trigger='netdev'
set system.wifi_led.dev='wlan0'
set system.wifi_led.mode='link tx rx'
set system.wifi_led.default='0'
commit system
__EOF

# add mac address from U-Boot partition to lan and wan devices
	MTD=`grep -e 'u-boot' /proc/mtd`
	MTD=`echo ${MTD} | sed 's/[a-z]*\([0-9]*\):.*/\1/'`
	[ -n "${MTD}" ] && {
		MACADDR=`dd if=/dev/mtdblock${MTD} bs=1 skip=262048 count=6 2>/dev/null | hexdump -e '1/1 "%02x"'`
		MACADDR2=$(( 0x${MACADDR} + 1))
		MACADDR2=`printf "%012x" ${MACADDR2}`

		MACADDR=`echo ${MACADDR} | sed 's/\(..\)/\1:/g' | sed 's/:$//'`
		MACADDR2=`echo ${MACADDR2} | sed 's/\(..\)/\1:/g' | sed 's/:$//'`

		uci set network.eth0.macaddr=${MACADDR}
		uci set network.lan.macaddr=${MACADDR}
		uci set network.wan.macaddr=${MACADDR2}
		uci commit network
	}
}

case "$(board_name)" in
 'Linksys WRT350N v2')
	wrt350nv2_default
	;;
esac
