# target/linux/mvebu/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

start() {
	. /lib/functions.sh

	case $(board_name) in
	linksys,caiman |\
	linksys,cobra |\
	linksys,mamba |\
	linksys,rango |\
	linksys,shelby |\
	linksys,venom)
		mtd resetbc s_env || true
		;;
	esac
}
