# target/linux/mvebu/base-files/etc/init.d/linksys_recovery
#!/bin/sh /etc/rc.common
# Copyright (C) 2015-2016 OpenWrt.org
# Copyright (C) 2016 LEDE-Project.org

START=97
boot() {
. /lib/functions.sh

case $(board_name) in
	armada-385-linksys-caiman|armada-385-linksys-cobra|armada-385-linksys-rango|armada-385-linksys-shelby|armada-385-linksys-venom|armada-xp-linksys-mamba)
		# make sure auto_recovery in uboot is always on
		AUTO_RECOVERY_ENA="`fw_printenv -n auto_recovery`"
		if [ "$AUTO_RECOVERY_ENA" != "yes" ] ; then
			fw_setenv auto_recovery yes
		fi
		# reset the boot counter
		mtd resetbc s_env
		;;
esac
}
