# customfeeds/packages/net/geth/files/geth.init
#!/bin/sh /etc/rc.common

START=90
STOP=10

USE_PROCD=1
PROG=/usr/bin/geth

start_service() {
    procd_open_instance
    procd_set_param command ${PROG}
    procd_append_param command --syncmode "light" --cache 1024
    procd_set_param respawn
    procd_close_instance
}
