# package/lean/wsdd2/files/wsdd2.init
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

start_service() {
    procd_open_instance
    procd_set_param command /usr/bin/wsdd2 -w
    procd_set_param respawn
    procd_set_param file /var/etc/smb.conf
    procd_close_instance
}

service_triggers() {
	PROCD_RELOAD_DELAY=3000
	procd_add_reload_trigger "dhcp" "system" "samba4"
}
