# customfeeds/packages/ipv6/miredo/files/miredo.init
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=41

PROG=/usr/sbin/miredo
CONFFILE=/etc/miredo/miredo.conf

start_service() {
  procd_open_instance
  procd_set_param command $PROG -f -c $CONFFILE
  procd_set_param respawn
  procd_set_param file $CONFFILE
  procd_set_param stdout 1 # forward stdout of the command to logd
  procd_set_param stderr 1 # same for stderr
  procd_close_instance
}
