# target/linux/ath79/mikrotik/base-files/etc/uci-defaults/04_led_migration
. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)
boardonly="${board##*,}"

case "$board" in
mikrotik,*)
	migrate_leds "^rb:=mikrotik:"
	;;
esac

migrations_apply system

exit 0
