# package/Twinzo/luci-app-chinadns/root/etc/init.d/chinadns
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2015 OpenWrt-dist
# Copyright (C) 2015 Jian Chang <aa65535@live.com>
#
# This is free software, licensed under the GNU General Public License v3.
# See /LICENSE for more information.
#

START=98
STOP=15

SERVICE_USE_PID=1
SERVICE_WRITE_PID=1
SERVICE_DAEMONIZE=1

port=$(uci get chinadns.@chinadns[0].port)
enabled=$(uci get chinadns.@chinadns[0].enable)
automod=$(uci get chinadns.@chinadns[0].automod)
autorefresh=$(uci get chinadns.@chinadns[0].autorefresh)
refreshtime=$(uci get chinadns.@chinadns[0].refreshtime)
boot_time=$(uci get chinadns.@chinadns[0].boot_time)

xappend() {
	parms="$parms $1"
}

append_bool() {
	local section="$1"
	local option="$2"
	local value="$3"
	local _loctmp
	config_get_bool _loctmp "$section" "$option"
	[ "$_loctmp" = 1 ] || return 0
	xappend "$value"
}

append_parm() {
	local section="$1"
	local option="$2"
	local switch="$3"
	local _loctmp
	config_get _loctmp "$section" "$option"
	[ -z "$_loctmp" ] && return 0
	xappend "$switch $_loctmp"
}

start_chinadns() {
	[ "$enabled" = 1 ] || return 0
	[ "$automod" = 1 ] && {
		uci del dhcp.@dnsmasq[-1].server
		uci del dhcp.@dnsmasq[-1].resolvfile
		uci set dhcp.@dnsmasq[-1].noresolv=1
		uci add_list dhcp.@dnsmasq[-1].server=127.0.0.1#$port
		uci commit dhcp
	}
	append_parm $1 port "-p"
	append_parm $1 server "-s"
	append_parm $1 chnroute "-c"
	append_bool $1 bidirectional "-d"
	service_start /usr/bin/chinadns \
		$parms -m
	return $?
}

boot() {
	[ $enabled = 0 ] && exit 0
	sleep "$boot_time"
	start
}

start() {
	config_load chinadns
	config_foreach start_chinadns chinadns
	if [ "$enabled" = "1" ]; then
		if [ "$autorefresh" = "1" ] ; then
			[ ! "$(cat /etc/crontabs/root |grep 'chinadns reload')" ] && echo "*/$refreshtime * * * * /etc/init.d/chinadns reload >/dev/null 2>&1" >> /etc/crontabs/root
		else
			sed -i '/chinadns reload/d' /etc/crontabs/root
		fi
	else
		sed -i '/chinadns reload/d' /etc/crontabs/root
	fi
	/etc/init.d/dnsmasq restart
}

stop() {
	[ "$enabled" = '0' ] && {
		[ "$automod" = '1' ] && {
			uci del dhcp.@dnsmasq[-1].server
			uci del dhcp.@dnsmasq[-1].noresolv
			uci set dhcp.@dnsmasq[-1].resolvfile=/tmp/resolv.conf.auto
			uci commit dhcp
		}	
	service_stop /usr/bin/chinadns
	/etc/init.d/dnsmasq restart
	}
}

