# customfeeds/packages/net/nsd/files/nsd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 OpenWrt.org

START=60
USE_PROCD=1
PROG=/usr/sbin/nsd

start_service() {
	mkdir -p /var/db/nsd
	chown network /var/db/nsd
	chmod 640 /etc/nsd/*.conf 2>/dev/null
	chgrp network /etc/nsd/*.conf 2>/dev/null
	procd_open_instance
	procd_set_param command "$PROG" -d
	procd_close_instance
}
