# package/lean/luci-app-nps/root/etc/uci-defaults/luci-nps
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@nps[-1]
	add ucitrack nps
	set ucitrack.@nps[-1].init=nps
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
