# customfeeds/luci/applications/luci-app-dynapoint/root/etc/uci-defaults/40_luci-dynapoint
#!/bin/sh

# needed for "Save and Apply" to restart dynapoint
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@dynapoint[-1]
	add ucitrack dynapoint
	set ucitrack.@dynapoint[-1].init="dynapoint"
	commit dynapoint
EOF

rm -f /tmp/luci-indexcache
exit 0

