# customfeeds/packages/multimedia/lcdgrilo/files/lcdgrilo.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009-2012 OpenWrt.org

START=60

start() {
	/usr/bin/lcdgrilo &
	echo $! > /var/run/lcdgrilo.pid
}

stop() {
	if [ ! -f /var/run/lcdgrilo.pid ]; then
		return
	fi
	/bin/kill `cat /var/run/lcdgrilo.pid`
}
