# target/linux/bcm27xx/base-files/etc/uci-defaults/50_rpi-4b-set-wifi
. /lib/functions.sh

case "$(board_name)" in
	raspberrypi,4-model-b)
		sed -i 's/VHT80/VHT40/g' /lib/wifi/mac80211.sh
		;;
esac

exit 0
