# target/linux/ipq806x/base-files/etc/hotplug.d/firmware/11-ath10k-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"ath10k/pre-cal-pci-0000:01:00.0.bin")
	case $board in
	buffalo,wxr-2533dhp)
		caldata_extract "ART" 0x1000 0x2f20
		ath10k_patch_mac $(mtd_get_mac_binary ART 0x1e)
		;;
	linksys,ea8500)
		caldata_extract "art" 0x1000 0x2f20
		ath10k_patch_mac $(macaddr_add $(mtd_get_mac_ascii devinfo hw_mac_addr) +1)
		;;
	nec,wg2600hp)
		caldata_extract "ART" 0x1000 0x2f20
		ath10k_patch_mac $(macaddr_add $(mtd_get_mac_binary PRODUCTDATA 0xc) +1)
		;;
	netgear,d7800 |\
	netgear,r7500v2 |\
	netgear,r7800)
		caldata_extract "art" 0x1000 0x2f20
		ath10k_patch_mac $(macaddr_add $(mtd_get_mac_binary art 0x6) +1)
		;;
	tplink,c2600)
		caldata_extract "radio" 0x1000 0x2f20
		ath10k_patch_mac $(macaddr_add $(mtd_get_mac_binary default-mac 0x8) -1)
		;;
	tplink,vr2600v)
		caldata_extract "ART" 0x1000 0x2f20
		ath10k_patch_mac $(macaddr_add $(mtd_get_mac_binary default-mac 0x0) -1)
		;;
	zyxel,nbg6817)
		caldata_extract "0:ART" 0x1000 0x2f20
		ath10k_patch_mac $(macaddr_add $(mtd_get_mac_ascii 0:APPSBLENV ethaddr) +1)
		;;
	esac
	;;
"ath10k/pre-cal-pci-0001:01:00.0.bin")
	case $board in
	buffalo,wxr-2533dhp)
		caldata_extract "ART" 0x5000 0x2f20
		ath10k_patch_mac $(mtd_get_mac_binary ART 0x18)
		;;
	linksys,ea8500)
		caldata_extract "art" 0x5000 0x2f20
		ath10k_patch_mac $(macaddr_add $(mtd_get_mac_ascii devinfo hw_mac_addr) +2)
		;;
	nec,wg2600hp)
		caldata_extract "ART" 0x5000 0x2f20
		ath10k_patch_mac $(mtd_get_mac_binary PRODUCTDATA 0xc)
		;;
	netgear,d7800 |\
	netgear,r7500v2 |\
	netgear,r7800)
		caldata_extract "art" 0x5000 0x2f20
		ath10k_patch_mac $(macaddr_add $(mtd_get_mac_binary art 0x6) +2)
		;;
	tplink,c2600)
		caldata_extract "radio" 0x5000 0x2f20
		ath10k_patch_mac $(mtd_get_mac_binary default-mac 0x8)
		;;
	tplink,vr2600v)
		caldata_extract "ART" 0x5000 0x2f20
		ath10k_patch_mac $(mtd_get_mac_binary default-mac 0x0)
		;;
	zyxel,nbg6817)
		caldata_extract "0:ART" 0x5000 0x2f20
		ath10k_patch_mac $(mtd_get_mac_ascii 0:APPSBLENV ethaddr)
		;;
	esac
	;;
*)
	exit 1
	;;
esac
