# customfeeds/packages/utils/smartmontools/files/smartd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=95

USE_PROCD=1
PROG=/usr/sbin/smartd

start_service() {
	procd_open_instance
	procd_set_param command "$PROG"
	procd_append_param command -n -q never
	procd_close_instance
}

reload_service() {
	procd_send_signal smartd
}
