# target/linux/qualcommax/ipq50xx/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	linksys,mx2000|\
	linksys,mx5500)
		mtd resetbc s_env || true
	;;
	redmi,ax3000)
	fw_printenv >/dev/null || return 0
	local flag_ota_reboot="$(fw_printenv flag_ota_reboot | cut -d'=' -f2)"
	if [[ "$flag_ota_reboot" == 1 ]]; then
		local mtd="$(grep -oE 'ubi.mtd=[a-zA-Z0-9\-\_]*' /proc/cmdline | cut -d'=' -f2)"
		case $mtd in
			rootfs)
				fw_setenv flag_last_success 0
				fw_setenv flag_ota_reboot 0
				fw_setenv flag_boot_success 1
				;;
			rootfs_1)
				fw_setenv flag_last_success 1
				fw_setenv flag_ota_reboot 0
				fw_setenv flag_boot_success 1
				;;
			*)
				echo "Unable to determine UBIPART: ubi.mtd=$mtd"
				return 1
				;;
		esac
	fi
	;;
	esac
}
