# customfeeds/luci/applications/luci-app-asterisk/root/etc/uci-defaults/40_luci-asterisk
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@asterisk[-1]
	add ucitrack asterisk
	set ucitrack.@asterisk[-1].init=asterisk
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
