# package/lean/luci-app-vlmcsd/root/etc/uci-defaults/luci-app-vlmcsd
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@vlmcsd[-1]
	add ucitrack vlmcsd
	set ucitrack.@vlmcsd[-1].init=vlmcsd
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
