# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/usr/bin/lcd2usb_echo &

board=$(cat /tmp/sysinfo/board_name | cut -d , -f2)
if [ ! -e /etc/firstboot_${board} ]; then
    /root/setup.sh
    touch /etc/firstboot_${board}
fi
/bin/mount -a

#iptables -t nat -A POSTROUTING -o eth0 -j MASQUERADE
exit 0
