# target/linux/mediatek/base-files/etc/hotplug.d/iface/99-mtk-lro
[ ifup = "$ACTION" ] && {
	[ -n "$DEVICE" ] && {
		if [ "$INTERFACE" == "lan" ]; then
			if [ -f /usr/sbin/ethtool ]; then
				ifname=eth0
				lan_ip=`uci -q get network.lan.ipaddr`
				ethdrv=`ethtool -i $ifname | grep mtk_soc_eth`
				[ -n "$ethdrv" ] && {
					ethtool -N $ifname flow-type tcp4 dst-ip $lan_ip loc 0
				}
			fi
		fi
	}
}
