# customfeeds/packages/sound/forked-daapd/files/forked-daapd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 OpenWrt.org

START=99
BIN=/usr/sbin/forked-daapd
PID=/var/run/forked-daapd.pid
SSD=start-stop-daemon

start() {
        $SSD -p $PID -S -x $BIN -- -P $PID
}
	
stop() {
        $SSD -p $PID -K -s SIGINT
}
