# customfeeds/packages/utils/uvol/files/uvol.init
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1
NAME=uvol
PROG=/usr/sbin/uvol

start_service() {
	[ "${__BOOT_UVOL}" = "1" ] || return 0
	procd_open_instance "$NAME"
	procd_set_param command "$PROG" boot
	procd_close_instance
}

boot() {
	__BOOT_UVOL=1
	start
}
