# target/linux/qualcommax/ipq50xx/base-files/etc/hotplug.d/firmware/10-ath11k-board_id
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

set_board_id() {
	echo "$1" >/lib/firmware/$FIRMWARE
}

board=$(board_name)

case "$FIRMWARE" in
	"ath11k/IPQ5018/hw1.0/board_id-ahb-c000000.wifi")
		case "$board" in
			redmi,ax3000)
				case "$(strings /dev/mtdblock15 | grep wl_pa_type= | cut -d= -f2)" in
					epa) set_board_id 0x24 ;;
					ipa|*) set_board_id 0x10 ;;
				esac
				;;
		esac
		;;
	*)
		exit 1
		;;
esac
