# target/linux/qualcommax/ipq50xx/base-files/etc/hotplug.d/firmware/11-ath11k-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"ath11k/IPQ5018/hw1.0/cal-ahb-c000000.wifi.bin")
	case "$board" in
	cmcc,rax3000q|\
	glinet,gl-b3000)
		caldata_extract "0:ART" 0x1000 0x20000
		addr=$(mtd_get_mac_binary "0:ART" 0x6)
		ath11k_patch_mac $(macaddr_add $addr 3) 0
		ath11k_set_macflag
		;;
	jdcloud,re-cs-03)
		caldata_extract_mmc "0:ART" 0x1000 0x20000
		;;
	linksys,mx2000|\
	linksys,mx5500)
		caldata_extract "0:ART" 0x1000 0x20000
		label_mac=$(mtd_get_mac_ascii devinfo hw_mac_addr)
		ath11k_patch_mac $(macaddr_add $label_mac 1) 0
		ath11k_remove_regdomain
		ath11k_set_macflag
		;;
	esac
	;;
"ath11k/QCN6122/hw1.0/cal-ahb-b00a040.wifi1.bin")
	case "$board" in
	cmcc,rax3000q|\
	glinet,gl-b3000)
		caldata_extract "0:ART" 0x26800 0x20000
		addr=$(mtd_get_mac_binary "0:ART" 0x6)
		ath11k_patch_mac $(macaddr_add $addr 4) 0
		ath11k_set_macflag
		;;
	jdcloud,re-cs-03)
		caldata_extract_mmc "0:ART" 0x26800 0x20000
		;;
	linksys,mx2000)
		caldata_extract "0:ART" 0x26800 0x20000
		label_mac=$(mtd_get_mac_ascii devinfo hw_mac_addr)
		ath11k_patch_mac $(macaddr_add $label_mac 2) 0
		ath11k_remove_regdomain
		ath11k_set_macflag
		;;
	esac
	;;
"ath11k/QCN9074/hw1.0/cal-pci-0001:01:00.0.bin")
	case "$board" in
	linksys,mx5500)
		caldata_extract "0:ART" 0x26800 0x20000
		label_mac=$(mtd_get_mac_ascii devinfo hw_mac_addr)
		ath11k_patch_mac $(macaddr_add $label_mac 2) 0
		ath11k_remove_regdomain
		ath11k_set_macflag
		;;
	esac
	;;
*)
	exit 1
	;;
esac
