# customfeeds/luci/applications/luci-app-pagekitec/root/etc/uci-defaults/40_luci-pagekitec
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@pagekitec[-1]
	add ucitrack pagekitec
	set ucitrack.@pagekitec[-1].init=pagekitec
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
