# target/linux/ramips/mt7621/base-files/etc/uci-defaults/01_enable_packet_steering
uci -q get network.globals.packet_steering > /dev/null || {
	uci set network.globals='globals'
	uci set network.globals.packet_steering=1
	uci commit network
}
