# target/linux/ramips/mt7621/base-files/etc/uci-defaults/04_led_migration
#!/bin/sh

. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
netgear,wndr3700-v5)
	migrate_leds "^wndr3700v5:=wndr3700-v5:"
	;;
esac

migrations_apply system

exit 0
