# target/linux/ramips/mt7621/base-files/etc/uci-defaults/04_led_migration
. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)
boardonly="${board##*,}"

case "$board" in
mikrotik,routerboard-750gr3)
	migrate_leds "^rb750gr3:=$boardonly:"
	;;
mikrotik,routerboard-m11g)
	migrate_leds "^rbm11g:=$boardonly:"
	;;
mikrotik,routerboard-m33g)
	migrate_leds "^rbm33g:=$boardonly:"
	;;
netgear,wndr3700-v5)
	migrate_leds "^wndr3700v5:=$boardonly:"
	;;
esac

migrations_apply system

exit 0
