# target/linux/lantiq/base-files/etc/hotplug.d/firmware/11-ath10k-caldata
#!/bin/sh
# Based on ar71xx 11-ath10k-caldata and 10-rt2x00-eeprom

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

case "$FIRMWARE" in
"ath10k/cal-pci-0000:02:00.0.bin")
	board=$(board_name)
	case $board in
		bt,homehub-v5a)
			caldata_extract_ubi "caldata" 0x5000 0x844
			ath10k_patch_mac $(macaddr_add $(mtd_get_mac_binary_ubi caldata 0x110c) +3)
			;;
		*)
			caldata_die "board $board is not supported yet"
			;;
	esac
	;;
esac
