# target/linux/ipq40xx/base-files/etc/hotplug.d/ieee80211/10_fix_wifi_mac
#!/bin/ash

[ "$ACTION" = "add" ] || exit 0

PHYNBR=${DEVPATH##*/phy}

[ -n $PHYNBR ] || exit 0

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)

case "$board" in
	hiwifi,c526a)
		[ "$PHYNBR" -eq 0 ] && \
			mtd_get_mac_ascii bdinfo "Vfac_mac " > /sys${DEVPATH}/macaddress
		;;
	*)
		;;
esac
