# package/my/luci-app-oscam/root/etc/init.d/oscam
#!/bin/sh /etc/rc.common

START=99
OSCAM_PATH=/usr/bin

[ -f  /tmp/oscam ] && 
{
	OSCAM_PATH=/tmp
	[ -x  /tmp/oscam ] ||  chmod +x /tmp/oscam
}

start() {
	enable=$(uci get oscam.config.enabled 2>/dev/null)
	if [ $enable -eq 1 ]; then	
	echo "enable"
	[ -d /var/log/oscam ] || mkdir -p /var/log/oscam
	service_start ${OSCAM_PATH}/oscam -b -r 2 -u
	fi
}

stop() {
	enable=$(uci get oscam.config.enabled 2>/dev/null)
	if [ $enable -ne 1 ]; then
	echo "stop!"
	killall oscam	
	killall -q -9 oscam	
	fi
}

restart(){
	stop
	sleep 2
	start
	echo "oscam Server has restarted."
}