# target/linux/mvebu/base-files/etc/uci-defaults/03_wireless
#!/bin/sh
#
# Copyright (C) 2015 OpenWrt.org
#

[ ! -e /etc/config/wireless ] && exit 0

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)

case "$board" in
armada-385-linksys-caiman|armada-385-linksys-cobra|armada-385-linksys-shelby|armada-xp-linksys-mamba)
	SKU=$(strings /dev/mtd3|sed -ne 's/^cert_region=//p')
	WIFIMAC2G=$(macaddr_add $(cat /sys/class/net/eth0/address) +1)
	WIFIMAC5G=$(macaddr_add $WIFIMAC2G +1)
	case "$SKU" in
		AP)
			REGD=CN
		;;
		AU)
			REGD=AU
		;;
		CA)
			REGD=CA
		;;
		EU)
			REGD=DE
		;;
		US)
			REGD=US
		;;
	esac

	case "$board" in
		armada-xp-linksys-mamba)
			WIFIMAC0=$WIFIMAC2G
			WIFIMAC1=$WIFIMAC5G
		;;
		*)
			WIFIMAC0=$WIFIMAC5G
			WIFIMAC1=$WIFIMAC2G
		;;
	esac

	uci get wireless.radio0.country || uci set wireless.radio0.country=$REGD
	uci get wireless.@wifi-iface[0].macaddr || uci set wireless.@wifi-iface[0].macaddr=$WIFIMAC0
	uci get wireless.radio1.country || uci set wireless.radio1.country=$REGD
	uci get wireless.@wifi-iface[1].macaddr || uci set wireless.@wifi-iface[1].macaddr=$WIFIMAC1
	;;
esac

uci commit wireless

exit 0
