# target/linux/uml/base-files/etc/inittab
# Copyright (c) 2013 The Linux Foundation. All rights reserved.
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
tty0::askfirst:/usr/libexec/login.sh
