# package/lean/luci-app-arpbind/root/etc/uci-defaults/luci-arpbind
#!/bin/sh
touch /etc/config/arpbind

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@arpbind[-1]
	add ucitrack arpbind
	set ucitrack.@arpbind[-1].init=arpbind
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
