# package/my/luci-app-macvlan/files/root/etc/uci-defaults/40_luci-app-macvlan
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@macvlan[-1]
	add ucitrack macvlan
	set ucitrack.@macvlan[-1].init=macvlan
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
