# package/diy/luci-app-fan/root/etc/init.d/luci-fan
#!/bin/sh /etc/rc.common

START=20

boot() {
    /usr/libexec/fan-control init
    start
}

start_instance() {
    local enabled 
    config_get_bool enabled $1 enabled 0
    [[ "$enabled" = 1 ]] || return 0

    config_get on_temp $1 on_temp
    config_get off_temp $1 off_temp

    [[ -n "$on_temp" ]] && on_temp=${on_temp}000
    [[ -n "$off_temp" ]] && off_temp=${off_temp}000

    /usr/libexec/fan-control set "$on_temp" $off_temp
}

start() {
    config_load luci-fan
    config_foreach start_instance luci-fan
    return 0
}
