# customfeeds/packages/net/foolsm/files/foolsm.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2010-2011 OpenWrt.org

START=45

SERVICE_USE_PID=1
SERVICE_PID_FILE=/var/run/foolsm.pid

start() {
	service_start /usr/sbin/foolsm -c /etc/foolsm/foolsm.conf -p $SERVICE_PID_FILE
}

stop() {
	service_stop /usr/sbin/foolsm
}

reload() {
	service_reload /usr/sbin/foolsm
}
