# customfeeds/packages/utils/telldus-core/files/telldusd
#!/bin/sh /etc/rc.common


START=99

USE_PROCD=1
PROG=/usr/sbin/telldusd
CONFFILE=/etc/tellstick.conf

start_service() {
	procd_open_instance
	procd_set_param command $PROG
	procd_append_param command --nodaemon # foreground required by procd
	procd_set_param file $CONFFILE
	procd_set_param term_timeout 10 # telldusd is slow to exit of some reason
	procd_set_param respawn  # respawn the service if it exits
	procd_set_param stdout 1 # forward stdout of the command to logd
	procd_set_param stderr 1 # same for stderr
	procd_close_instance
}
