# package/lean/luci-theme-amazeui/root/etc/uci-defaults/luci-theme-amazeui
#!/bin/sh
uci batch <<-EOF
	set luci.themes.amazeui=/luci-static/amazeui
        commit luci
EOF
	
