# customfeeds/luci/applications/luci-app-polipo/root/etc/uci-defaults/40_luci-polipo
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@polipo[-1]
	add ucitrack polipo
	set ucitrack.@polipo[-1].init=polipo
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
