# target/linux/mediatek/filogic/base-files/etc/hotplug.d/firmware/11-mt76-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"mediatek/mt7986_eeprom_mt7976_dbdc.bin")
	case "$board" in
	asus,tuf-ax4200)
		CI_UBIPART="UBI_DEV"
		caldata_extract_ubi "Factory" 0x0 0x1000
		;;
	esac
	;;
"mediatek/mt7986_eeprom_mt7976_dual.bin")
	case "$board" in
	glinet,gl-mt6000)
		caldata_extract_mmc "factory" 0x0 0x1000
		;;
	esac
	;;
*)
	exit 1
	;;
esac
