# customfeeds/packages/utils/haveged/files/haveged.init
#!/bin/sh /etc/rc.common

START=01
USE_PROCD=1

HAVEGED_THRESHOLD=1024
HAVEGED_DCACHE=32
HAVEGED_ICACHE=32

start_service() {
	procd_open_instance
	procd_set_param command /usr/sbin/haveged
	procd_append_param command -F -w $HAVEGED_THRESHOLD -d $HAVEGED_DCACHE -i $HAVEGED_ICACHE -v 1
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}
