# customfeeds/packages/admin/schroot/files/schroot.init
#!/bin/sh /etc/rc.common

START=94

start() {
	[ -d /var/lib/schroot ] || {
		mkdir -p /var/lib/schroot/mount /var/lib/schroot/session /var/lib/schroot/unpack /var/lib/schroot/union/overlay /var/lib/schroot/union/underlay
		chmod 0700 /var/lib/schroot
	}
}
