# package/Twinzo/luci-app-oafControl/root/etc/uci-defaults/luci-oafControl
#!/bin/sh

mkdir /tmp/oafControl 2>/dev/null
mkdir /etc/oafControl 2>/dev/null
rm -f /tmp/luci-indexcache
chmod 755 /etc/init.d/oafControl
chmod 755 /usr/bin/oafControl.sh
exit 0