# package/Twinzo/luci-app-oafControl/root/etc/uci-defaults/luci-oafControl
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@oafControl[-1]
	add ucitrack oafControl
	set ucitrack.@oafControl[-1].init=oafControl
	commit ucitrack
EOF

mkdir /tmp/oafControl 2>/dev/null
mkdir /etc/oafControl 2>/dev/null
rm -f /tmp/luci-indexcache
chmod 755 /etc/init.d/oafControl
chmod 755 /usr/bin/oafControl.sh
exit 0