# package/luci-app-netspeedtest/root/etc/uci-defaults/40_luci-netspeedtest
#!/bin/sh
chmod +x /etc/init.d/netspeedtest  >/dev/null 2>&1
[ -f "/etc/config/netspeedtest" ] || {
cat >/etc/config/netspeedtest <<-EOF
config netspeedtest 'netspeedtest'
        option enabled '0'
	option port '3300'

config homebox 'homebox'
        option enabled '0'
	option port '3300'

config speedtestiperf3  'speedtestiperf3'

config speedtestwan 'speedtestwan'
	
config speedtestport 'speedtestport'
	option sport '443'
	option sdomain 'www.taobao.com'
EOF
}
[  `uci -q get netspeedtest.global` ]  && uci set netspeedtest.global=global
[  `uci -q get netspeedtest.netspeedtest` ]  || uci set netspeedtest.netspeedtest=netspeedtest
[  `uci -q get netspeedtest.speedtestiperf3` ] || uci set netspeedtest.speedtestiperf3=speedtestiperf3
[  `uci -q get netspeedtest.speedtestport` ]  || uci set netspeedtest.speedtestport=speedtestport
[  `uci -q get netspeedtest.speedtestwan` ]  || uci set netspeedtest.speedtestport=speedtestwan
LOGD=/etc/netspeedtest
LOG=/var/log/netspeedtest.log
LOGE=$LOGD/netspeedtest.log
LOGT=$LOGD/netspeedtestpos
[ -d "$LOGD" ] || mkdir -p $LOGD
[ -f "$LOGE" ] || echo "start" > $LOGE
[ -f "$LOGT" ] || echo "1" > $LOGT
[ -f "$LOG" ] || echo "start" > $LOG

rm -rf /tmp/luci-modulecache /tmp/luci-indexcache*
exit 0
