# customfeeds/luci/applications/luci-app-splash/root/etc/hotplug.d/iface/25-splash
#!/bin/sh

[ "$ACTION" = ifup ] || exit 0

/etc/init.d/firewall enabled || exit 0

if [ -x /etc/init.d/luci_splash ]; then
	restart_splash() {
	        local net="$1"
		if [ "$INTERFACE" = "$net" ]; then
			logger -t splash "Reloading splash firewall rules due to ifup of $INTERFACE ($DEVICE)"
			/etc/init.d/luci_splash restart
		fi
	}
	config_load luci_splash
	config_foreach restart_splash iface
fi
