# customfeeds/packages/utils/pcmciautils/files/etc/init.d/pcmcia_socket
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org 

START=41
start() {
	if [ -d /sys/class/pcmcia_socket/ ]; then
		# 2.6 kernels
		cd /sys/class/pcmcia_socket/
		for PCMCIA_SOCKET in *; do
			if [ $PCMCIA_SOCKET == \* ]; then
				continue;
			fi;
			SOCKET_NO=$(echo $PCMCIA_SOCKET | cut -b14-)
			ACTION=add
			export SOCKET_NO ACTION
			/sbin/hotplug-call pcmcia_socket
		done;
	fi;
}
