# target/linux/meson/base-files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

bash /root/resize.sh
exit 0
