# package/lean/dogcom/files/etc/init.d/dogcom_log
#!/bin/sh /etc/rc.common
#Cpyright by Twizo<1282055288@qq.com>
START=90

stochastic(){
	echo `dd if=/dev/urandom bs=1 count=32 2>/dev/null | md5sum | cut -b 0-12 | sed 's/\(..\)/\1:/g; s/.$//'`
}

minus_mac()
{
	d_mac=$(uci get dogcom.@dogcom[0].mac 2>/dev/null)
	l_mac=0x$(echo "$d_mac"|sed 's/\(.*\)\(.\)$/\2/')
	dec2hex(){
     		printf "%x" $1
 		}
##十进制转十六进制
	l_mac=`echo $((0xb^0x1))` ##异或求反
	mac=$(echo $d_mac|sed "s/.$/$(dec2hex $l_mac)/")
	echo $mac
}


sto(){ 
#	en_JumpTwoMac=$(uci get dogcom.@dogcom[0].en_JumpTwoMac 2>/dev/null)
	local stom
#	if [ "$en_JumpTwoMac" -eq "1" ]; then
#		stom=$("minus_mac")
#	else
#		stom=$("stochastic")
#	fi
	uci set network.wan.macaddr="$stom"
	uci set network.wan_dev.macaddr="$stom" 2>/dev/null
#	uci set network.wan6.macaddr="$stom"
	uci commit
}
#mac只会在两个地址之中变化

start()
{	
	netcheck = $(uci get dogcom.@dogcom[0].netcheck 2>/dev/null)
	if [ "$netcheck"x != "1"x ]; then
		return
	fi
	
	tries=0
	echo --- my_watchdog start ---
	while [[ $tries -lt 2 ]]
	do
	        if /bin/ping -c 1 8.8.8.8 >/dev/null
		then
			echo "网络连接正常，持续监测中" >> /tmp/dogcom.log
	        	echo --- exit ---
			exit 0
       	 	fi
        	tries=$((tries+1))
        	sleep 3
	done
	echo "网络连接不正常，开始修改mac地址并重新连接" >> /tmp/dogcom.log
	sto
	/etc/init.d/network reload

}
