# target/linux/ramips/mt76x8/base-files/etc/uci-defaults/04_led_migration
#!/bin/sh

. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
tplink,archer-c20-v4)
	migrate_leds "^c20-v4:=archer-c20-v4:"
	;;
tplink,archer-c50-v3)
	migrate_leds "^c50-v3:=archer-c50-v3:"
	;;
tplink,archer-c50-v4)
	migrate_leds "^c50-v4:=archer-c50-v4:"
	;;
esac

migrations_apply system

exit 0
