# customfeeds/luci/applications/luci-app-udpxy/root/etc/uci-defaults/40_luci-udpxy
#!/bin/sh

# needed for "Save and Apply" to restart udpxy
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@udpxy[-1]
	add ucitrack udpxy
	set ucitrack.@udpxy[-1].init="udpxy"
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
