# customfeeds/luci/applications/luci-app-dnscrypt-proxy/root/etc/uci-defaults/60_luci-dnscrypt-proxy
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@dnscrypt-proxy[-1]
	add ucitrack dnscrypt-proxy
	set ucitrack.@dnscrypt-proxy[-1].init=dnscrypt-proxy
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
