# customfeeds/packages/admin/monit/files/monit.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 OpenWrt.org

START=60
USE_PROCD=1
PROG=/usr/bin/monit

start_service() {
	[ -f /etc/monitrc ] || { echo "monit: /etc/monitrc is missing"; return 1; }
	chmod 0600 /etc/monitrc
	procd_open_instance
	# -I runs in foreground, as procd requires
	procd_set_param command "$PROG" -I
	procd_close_instance
}
