# customfeeds/luci/modules/luci-newapi/root/etc/uci-defaults/099-newapi-copy
#!/bin/sh

cp -a /etc/newapi/* /

exit 0
