# LuCI-Rosy-master/root/etc/uci-defaults/30_luci-theme-rosy
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Rosy=/luci-static/rosy
	set luci.main.mediaurlbase=/luci-static/rosy
	commit luci
EOF
exit 0
