# package/Twinzo/luci-app-idns/root/etc/init.d/idns
#!/bin/sh /etc/rc.common
#source /lib/functions.sh
#----------------------------------------------
#Copyright (C) 2019 Twinzo1 <1282055288@qq.com>
#----------------------------------------------
START=98

#主路由或开启dhcp的路由器
set_link(){
	link_leases=`uci get idns.@main_router[0].leases 2>/dev/null`
	comments_file=`uci get idns.@main_router[0].comment_leases 2>/dev/null`
	ln -s /tmp/dhcp.leases /www/$link_leases 2>/dev/null
	[ -n "$comment_leases" ] && ln -s /tmp/dhcp_comments.conf /www/$comment_leases 2>/dev/null

}

#中继或未开启dhcp路由器
get_link(){
	upstreamIP=`uci get idns.@relay_router[0].upstreamIP`
	leases_name=`uci get idns.@relay_router[0].leases`
	comments_file=`uci get idns.@relay_router[0].comment_leases 2>/dev/null`
	curl $upstreamIP/$leases_name > /tmp/dhcp.leases
	[ -n "$comments_file" ] && curl $upstreamIP/$comments_file > /tmp/dhcp_comments.conf
}

redirect(){
	config_redir(){
		config_get enable $1 enabled
		config_get ipaddr $1 ipaddr
		config_get redirIP $1 web
		if [ "$enable" = "1" ]; then
			echo "iptables -t nat -A PREROUTING -d $ipaddr -p tcp -m tcp --dport 80 -j DNAT --to-destination $redirIP	#IDNS REDIRCTION$ipaddr" >> /etc/firewall.user
		fi
	}
	sed -i "/#IDNS REDIRCTION/d" /etc/firewall.user 2 >/dev/null
	config_load idns
	config_foreach config_redir redirection
	/etc/init.d/firewall restart
}

init_redir(){
	config_redir(){
		config_get enable $1 enabled
		[ $enable = 1 ] && echo 1
	}
	config_load idns
	config_foreach config_redir redirection
}

start(){
	IDNS=$(uci get idns.@idns[0].domain )
	main_enable=`uci get idns.@main_router[0].enabled`
	relay_enable=`uci get idns.@relay_router[0].enabled`
	
	[ -z $IDNS -a $main_enable !=1 -a $relay_enable != 1 -a -z $(init_redir) ] && /etc/init.d/idns disable && sed -i '/#added by IDNS/d' /etc/hosts && return
	[ $main_enable = 1 ] && set_link
	[ $relay_enable = 1 ] && get_link
	[ -n $(init_redir) ] && redirect
	if [ -n "$IDNS" ]; then
		lan_ipaddr=$(uci get network.lan.ipaddr)
		sed -i '/#added by IDNS/d' /etc/hosts
		echo  "$lan_ipaddr $IDNS #added by IDNS " >> /etc/hosts
		/etc/init.d/dnsmasq restart
	else
		isNull=`grep "added by IDNS" /etc/firewall.user`
		if [ -n "$isNull" ]; then
			sed -i '/#added by IDNS/d' /etc/hosts
			/etc/init.d/dnsmasq restart
		fi
	fi
	/etc/init.d/idns enable
}
