# customfeeds/packages/mail/dovecot/files/dovecot.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2013 OpenWrt.org

START=75
STOP=75

USE_PROCD=1

start_service() {
	mkdir -p -m 0755 /var/lib/dovecot
	mkdir -p -m 0755 /var/run/dovecot
	chmod 0750 /var/lib/dovecot
	procd_open_instance
	procd_set_param command /usr/sbin/dovecot -F
	procd_close_instance
}
