# customfeeds/packages/utils/uledd/files/uledd.init
#!/bin/sh /etc/rc.common

START=95

USE_PROCD=1
PROG=/sbin/uledd

start_service() {
	procd_open_instance
	procd_set_param command "$PROG"
	procd_set_param respawn
	procd_close_instance
}
