# customfeeds/packages/net/ngircd/files/ngircd.init
#!/bin/sh /etc/rc.common

START=65

start() {
	service_start /usr/sbin/ngircd
}

stop() {
	service_stop /usr/sbin/ngircd
}
