# customfeeds/luci/applications/luci-app-shadowsocks-libev/root/etc/uci-defaults/40_luci-shadowsocks-libev
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@shadowsocks-libev[-1]
	add ucitrack shadowsocks-libev
	set ucitrack.@shadowsocks-libev[-1].init=shadowsocks-libev
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
