# package/lean/luci-app-mwan3helper/root/etc/uci-defaults/luci-mwan3helper
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@mwan3helper[-1]
	add ucitrack mwan3helper
	set ucitrack.@mwan3helper[-1].init=mwan3helper
	commit ucitrack
EOF

rm -f /tmp/luci*
exit 0
