# package/lienol/luci-app-control-mia/root/etc/uci-defaults/luci-app-control-mia
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@mia[-1]
	add ucitrack mia
	set ucitrack.@mia[-1].init=mia
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
