# package/base-files/files/etc/hotplug.d/leds/00-init
#!/bin/sh

[ "$ACTION" = "add" ] && /etc/init.d/led start "$DEVICENAME"
