# package/lean/luci-app-flowoffload/root/etc/init.d/flowoffload
#!/bin/sh /etc/rc.common
# Copyright (c) 2011-2015 OpenWrt.org

START=99

restart(){
	/etc/init.d/firewall restart
}


