# package/lean/luci-app-adbyby-plus/root/etc/uci-defaults/adbyby
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@adbyby[-1]
	add ucitrack adbyby
	set ucitrack.@adbyby[-1].init=adbyby
	commit ucitrack
	delete firewall.adbyby
	set firewall.adbyby=include
	set firewall.adbyby.type=script
	set firewall.adbyby.path=/usr/share/adbyby/firewall.include
	set firewall.adbyby.reload=1
	commit firewall
EOF

mkdir -p /etc/dnsmasq.d

/etc/init.d/adbyby stop
/etc/init.d/adbyby enable

rm -f /tmp/luci-indexcache
exit 0
